# sys-apps/dbus/files/dbus
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/sys-apps/dbus/files/dbus,v 1.6 2005/11/08 04:46:00 compnerd Exp $

start() {

	ebegin "Starting D-BUS system messagebus"

	start-stop-daemon --start --pidfile /var/run/dbus.pid --exec /usr/bin/dbus-daemon-1 -- --system

	eend $?

}

stop() {

	local retval

	ebegin "Stopping D-BUS system messagebus"

	start-stop-daemon --stop --pidfile /var/run/dbus.pid
	retval=$?

	eend ${retval}

	# This should vanish when baselayout-1.12.0 is marked stable.
	[[ -f /var/run/dbus.pid ]] && rm -f /var/run/dbus.pid

   # Changes to the system dbus socket run file force us to take care of both
   [[ -f /var/lib/dbus/system_bus_socket ]] && rm -rf /var/lib/dbus/system_bus_socket
   [[ -f /var/run/dbus/system_bus_socket ]] && rm -rf /var/run/dbus/system_bus_socket

	return ${retval}

}
