# sys-apps/dbus/files/dbus
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/sys-apps/dbus/files/dbus,v 1.5 2005/03/11 21:21:09 azarah Exp $

start() {

	ebegin "Starting D-BUS system messagebus"

	start-stop-daemon --start --pidfile /var/run/dbus.pid --exec /usr/bin/dbus-daemon-1 -- --system

	eend $?

}

stop() {

	local retval

	ebegin "Stopping D-BUS system messagebus"

	start-stop-daemon --stop --pidfile /var/run/dbus.pid
	retval=$?

	eend ${retval}

	# This should vanish when baselayout-1.12.0 is marked stable.
	rm -f /var/run/dbus.pid

	return ${retval}

}
