# dev-util/cvsd/files/cvsd.rc6
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-util/cvsd/files/cvsd.rc6,v 1.2 2004/03/04 18:40:19 vapier Exp $

depend() {
	need net
}

start() {
	ebegin "Starting cvsd"
	start-stop-daemon --start --quiet --background --exec /usr/sbin/cvsd
	eend $?
}

stop() {
	ebegin "Stopping cvsd"
	start-stop-daemon --stop --quiet --pidfile /var/run/cvsd.pid
	eend $?
}
