# net-misc/radvd/files/radvd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/radvd/files/radvd,v 1.4 2004/07/15 00:10:29 agriffis Exp $

depend () {
	need net
}

start () {
	ebegin "Starting radvd"
	echo 1 > /proc/sys/net/ipv6/conf/all/forwarding
	start-stop-daemon --start --quiet --exec /usr/sbin/radvd -- -u radvd -p /var/run/radvd/radvd.pid
	eend
}

stop() {
	ebegin "Stopping radvd"
	echo 0 > /proc/sys/net/ipv6/conf/all/forwarding
	start-stop-daemon --stop --quiet --pidfile /var/run/radvd/radvd.pid
	eend
}
