# net-misc/radvd/files/radvd
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/radvd/files/radvd,v 1.1 2002/10/18 18:59:42 hannes Exp $

depend () {
	need net
}

start () {
	ebegin "Starting radvd"
	start-stop-daemon --start --quiet --exec /usr/sbin/radvd -- -u radvd -p /var/run/radvd/radvd.pid
	eend
}

stop() {
	ebegin "Stopping radvd"
	start-stop-daemon --stop --quiet --pidfile /var/run/radvd/radvd.pid
	eend
}
