# app-admin/bacula/files/bacula-init3
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/bacula/files/bacula-init3,v 1.1 2005/07/03 03:58:53 fserb Exp $

depend() {
		need net
		%%USE_DB%%
}

start() {
	local retval=0
	for s in ${SERVICES} ; do
		case ${s} in
			sd)
				ebegin "Starting bacula storage daemon"
				start-stop-daemon --start --quiet --exec /usr/sbin/bacula-sd \
					-- ${SD_OPTIONS}
				retval=$(( $retval + $? ))
				;;
			fd)
				ebegin "Starting bacula file daemon"
				start-stop-daemon --start --quiet --exec /usr/sbin/bacula-fd \
					-- ${FD_OPTIONS}
				retval=$(( $retval + $? ))
				;;
			dir)
				ebegin "Starting bacula director"
				start-stop-daemon --start --quiet --exec /usr/sbin/bacula-dir \
					-- ${DIR_OPTIONS}
				retval=$(( $retval + $? ))
				;;		
		esac
	done
	eend $retval
}

stop() {
	local retval=0
	for s in ${SERVICES} ; do
		case ${s} in
			sd)
				ebegin "Stopping bacula storage daemon"
				start-stop-daemon --stop --quiet --pidfile /var/run/bacula-sd.*.pid 
				retval=$(( $retval + $? ))
				;;
			fd)
				ebegin "Stopping bacula file daemon"
				start-stop-daemon --stop --quiet --pidfile /var/run/bacula-fd.*.pid
				retval=$(( $retval + $? ))
				;;
			dir)
				ebegin "Stopping bacula director"
				start-stop-daemon --stop --quiet --pidfile /var/run/bacula-dir.*.pid
				retval=$(( $retval + $? ))
				;;		
		esac
	done
	eend $retval
}
