# app-admin/msyslog/files/msyslog-init
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/msyslog/files/msyslog-init,v 1.8 2003/02/14 22:13:32 vapier Exp $

depend() {
	need clock hostname
	provide logger
}

start() {
	ebegin "Starting msyslog"
	start-stop-daemon --start --quiet --exec /usr/sbin/msyslogd \
		-- ${MSYSLOG_OPTS}
	eend $? "Failed to start msyslog"
}

stop() {
	ebegin "Stopping msyslog"
	start-stop-daemon --stop --quiet --pidfile /var/run/msyslogd.pid
	eend $? "Failed to stop msyslog"
}
