# net-p2p/fms/files/fms
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-p2p/fms/files/fms,v 1.1 2009/02/05 18:36:43 tommy Exp $

start() {
	ebegin "Starting FMS"
	if [ ! -f /etc/init.d/sysfs ]; then
		start-stop-daemon -S -u freenet:freenet -d /var/freenet/fms -x ./fms --daemon
		eend $?
	else
		cd /var/freenet/fms
		su -s /bin/sh -c "./fms --daemon" freenet
		eend $?
	fi
}

stop() {
	ebegin "Stopping FMS"
	if [ ! -f /etc/init.d/sysfs ]; then
		start-stop-daemon -K -u freenet:freenet -d /var/freenet/fms ./fms --daemon
		eend $?
	else
		pkill -u freenet -f "fms --daemon"
		eend 0
	fi
	if RC_CMD=restart; then
		ebegin "Waiting a few seconds to be sure FMS is ready for restart"
		sleep 10
		#give fms time to stop and release the nntp port
	fi
}

