# sys-libs/gpm/files/gpm-1.19.3-r4
#!/sbin/runscript
PIDFILE=/var/run/gpm.pid
EXE=/usr/sbin/gpm

start() {
		ebegin "Starting ${SERVICE}"
     		start-stop-daemon --start --quiet --exec ${EXE} -t ${MOUSE}
		eend ${?} "Started ${SERVICE}." "Error starting ${SERVICE}."
}

stop() {
		ebegin "Stopping ${SERVICE}"
		start-stop-daemon --stop --quiet --pidfile ${PIDFILE} --exec ${EXE}
		eend ${?} "Stopped ${SERVICE}." "Error stopping ${SERVICE}."
}

restart() {
		stop
		start
}
