# sys-apps/selinux-small/files/scmpd
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-apps/selinux-small/files/scmpd,v 1.1 2003/06/13 16:49:23 pebenito Exp $

depend() {
        need net
}

start () {
	ebegin "Staring scmpd"
	start-stop-daemon --start --quiet --pidfile /var/run/scmpd.pid \
				--startas /usr/sbin/scmpd
	eend $?
}

stop () {
	eend "Stopping scmpd"
	start-stop-daemon --stop --quiet --pidfile /var/run/scmpd.pid
	eend $?
}
