# sys-power/powerman/files/powerman
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-power/powerman/files/powerman,v 1.1 2008/09/12 18:14:52 weaver Exp $

depend() {
	need net
}

start() {
	ebegin "Starting powermand"
	start-stop-daemon --start --exec /usr/sbin/powermand
	eend $?
}

stop() {
	ebegin "Stopping powermand"
	start-stop-daemon --stop --exec /usr/sbin/powermand
	eend $?
}
