# sys-power/powerman/files/powerman
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-power/powerman/files/powerman,v 1.4 2014/05/12 10:19:16 ssuominen Exp $

depend() {
	need net
}

start() {
	ebegin "Starting powermand"

	checkpath -q -d -m 0755 /var/run/powerman

	start-stop-daemon --start \
		--pidfile /var/run/powerman/powermand.pid \
		--exec /usr/sbin/powermand
	eend $?
}

stop() {
	ebegin "Stopping powermand"
	start-stop-daemon --stop \
		--pidfile /var/run/powerman/powermand.pid \
		--exec /usr/sbin/powermand
	eend $?
}
