# mail-mta/qmail-ldap/files/1.03-r4/qmail-control
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/mail-mta/qmail-ldap/files/1.03-r4/qmail-control,v 1.2 2004/07/14 22:34:32 agriffis Exp $

opts="start stop restart reload"

start() {
    ebegin "Starting qmail mta ..."
    svc -u /service/qmail-{smtpd,send,qmqp,qmtp}
    eend $?
    ebegin "Starting qmail mta logging ..."
    svc -u /service/qmail-{smtpd,send,qmqp,qmtp}/log
    eend $?
}

stop() {
    ebegin "Stopping qmail mta ..."
    svc -d /service/qmail-{smtpd,send,qmqp,qmtp}
    eend $?
    ebegin "Stopping qmail mta logging ..."
    svc -d /service/qmail-{smtpd,send,qmqp,qmtp}/log
    eend $?
}

reload() {
    ebegin "Reloading 'locals' and 'virtualdomains' control files."
    svc -h /service/qmail-send
    eend $?
}
