# app-pda/synce-odccm/files/init.odccm
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	after bluetooth hald
}


start() {
   	ebegin "Starting Synce"
	ebegin "    Starting odccm"
	start-stop-daemon --start --quiet --exec /usr/sbin/odccm
	eend $?
	eend 0
}

stop() {
	ebegin "Shutting down Synce"
	ebegin "    Stopping odccm"
	start-stop-daemon --stop --quiet --exec /usr/sbin/odccm
	eend $?
	eend 0
}

# vim:ts=4
