# net-p2p/btg/files/btgd-init
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-p2p/btg/files/btgd-init,v 1.1 2007/10/14 23:37:14 angelos Exp $ 

BTG_HOME="$(getent passwd ${BTG_USER} | cut -d : -f 6)"

depend() {
	need net
}

checkconfig() {
	if [ -z ${BTG_USER} ]; then
		eerror "BTG_USER unset, please edit /etc/conf.d/btgd first"
		return 1
	elif [ ! -d "${BTG_HOME}" ]; then
		eerror "The home for ${BTG_USER} does not exist"
		eerror "Please check /etc/conf.d/btgd"
		return 1
	elif [ ! -f "${BTG_HOME}/.btg/daemon.ini" ]; then
		eerror "Unable to find daemon.ini in ${BTG_HOME}/.btg"
		eerror "Please run btg-config first"
		return 1
	fi
}

start() {
	ebegin "Starting BTG Daemon"
	checkconfig || return 1

	start-stop-daemon --start \
		-c ${BTG_USER} --quiet -b --make-pidfile \
		--pidfile /var/run/btgdaemon.pid \
		-x /usr/bin/btgdaemon -- -n
	eend $?
}

stop() {
	ebegin "Stopping BTG Daemon"
	start-stop-daemon --stop --quiet --pidfile /var/run/btgdaemon.pid
	eend $?
}
