# net-misc/cfengine/files/cf-monitord.rc6
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting cf-monitord"
	start-stop-daemon --start --exec /usr/sbin/cf-monitord -- ${EXTRA_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping cf-monitord"
	start-stop-daemon --stop --exec /usr/sbin/cf-monitord
	eend $?
}
