# mail-mta/sendmail/files/sendmail
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/mail-mta/sendmail/files/sendmail,v 1.2 2004/07/14 22:36:41 agriffis Exp $

depend() {
	need net
	use logger
	provide mta
}

start() {
	ebegin "Starting sendmail"
	/usr/bin/newaliases > /dev/null 2>&1
	(cd /var/spool/mqueue; rm -f xf*)
	/usr/sbin/sendmail ${SENDMAIL_OPTS} > /dev/null 2>&1
	/usr/sbin/sendmail ${CLIENTMQUEUE_OPTS} > /dev/null 2>&1
	eend $?
}

stop() {
	ebegin "Stopping sendmail"
	killall ${KILL_OPTS} sendmail
	eend $?
}
