# sys-apps/isapnptools/files/isapnp.rc6
#!/sbin/runscript

depend() {
	need localmount
}

start() {
if [ -e /etc/isapnp.conf ]; then
	if [ -e /sbin/isapnp ] ; then
		ebegin "Configuring ISA PnP devices"
		/sbin/isapnp /etc/isapnp.conf 1>&2
		eend $? "Error configuring ISA PnP"
	else
		eerror "You're missing /sbin/isapnp. Do emerge sys-apps/isapnptools"
		return 1
	fi
else
	eerror "You need an /etc/isapnp.conf file first"
	eerror "There is an example in docs/isapnptools/txt/isapnp.conf.gz"
	return 1
fi
}
