# sys-apps/microcode_ctl/files/1.06/microcode_ctl
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/sys-apps/microcode_ctl/files/1.06/microcode_ctl,v 1.1 2002/04/30 13:52:59 sandymac Exp $


depend() {
	need localmount
}

start() {
	ebegin "Updating Intel Microcode"
	if [ -f "/dev/cpu/microcode" ]
	then
		microcode_ctl -iu
		eend $? "Failed to update microcode"
	else
		eend 1 "/dev/cpu/microcode not found! Please check your kernel config."
	fi
}

# vim:ts=4
