# net-analyzer/hyperic-hq-agent/files/hyperic-hq-agent.init-2.0.22
#!/sbin/runscript

depend() {
	need net clock hostname
}

start() {
	ebegin "Starting Hyperic HQ Agent"
	/opt/hyperic-hq-agent-2.0.22/hq-agent.sh start
	eend $?
}

stop() {
	ebegin "Stopping Hyperic HQ Agent"
	/opt/hyperic-hq-agent-2.0.22/hq-agent.sh stop
	eend $?
}
