# sys-apps/at/files/atd.rc6
#!/sbin/runscript

depend() {
 	need clock
}

function start() {
	# Start daemons.
	ebegin "Starting atd..."
	start-stop-daemon --start --quiet --exec /usr/sbin/atd
	eend $? "Started atd." "Error Starting atd"
}

function stop() {
	# Stop daemons.
	ebegin "Shutting down atd..."
	start-stop-daemon --stop --quiet --pidfile /var/run/atd.pid
	eend $? "Stopped atd." "Error Stopping atd."
}
