# sys-apps/at/files/atd.rc6
#!/sbin/runscript

depend() {
    need clock
}

function start() {
    ebegin "Starting atd"
    start-stop-daemon --start --quiet --exec /usr/sbin/atd
    eend $?
}

function stop() {
    ebegin "Shutting down atd"
    start-stop-daemon --stop --quiet --pidfile /var/run/atd.pid
    eend $?
}
