# sys-apps/at/files/atd.rc6
#!/sbin/runscript

depend() {
    need clock
}

start() {
    ebegin "Starting atd"
    start-stop-daemon --start --quiet --exec /usr/sbin/atd
    eend $?
}

stop() {
    ebegin "Shutting down atd"
    #dont kill with --pidfile, as it does not always get created -- Azarah
    start-stop-daemon --stop --quiet --exec /usr/sbin/atd
    eend $?
}
