# sys-apps/at/files/atd.rc6
#!/sbin/runscript

SERVICE=atd
EXE=/usr/sbin/atd
PIDFILE=/var/run/atd.pid
opts="start stop restart"

# See how we were called.

start() {

	ebegin "Starting ${SERVICE}..."
	start-stop-daemon --start --quiet --exec ${EXE}
	eend $?

}

stop() {

	ebegin "Shutting down ${SERVICE}..."
	start-stop-daemon --stop --quiet --pidfile ${PIDFILE} 
	eend $?

}

restart () {

	stop
	start

}
