# net-analyzer/mwcollect/files/confd
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
#
# options for mwcollectd
#
# Usage: /usr/sbin/mwcollectd [options]
# Where options can be one of the following:
#
# --console-log [=tagpattern] - set console logging level
# --chroot=/opt/mwcollect/chroot - setup a chroot

OPTIONS=""

# see http://www.mwcollect.org/wiki/DocConfiguration for
# more info, which is currently empty, though
