# net-analyzer/mwcollect/files/initd
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later

opts="depend checkconfig start stop"

depend() {
	need net 
}

checkconfig() {
	if [ ! -e /etc/mwcollect/mwcollect.conf ] ; then
		eerror "/etc/mwcollect/mwcollect.conf not found!"
		return 1
	fi
}

start() {
	ebegin "Starting mwcollectd"

	checkconfig && \
	start-stop-daemon --start --quiet \
		--exec /usr/sbin/mwcollectd -- --pid-file=/var/run/mwcollect.pid \
		--daemon --user=nobody --config=/etc/mwcollect/mwcollect.conf \
		${OPTIONS}
	eend $?
}

stop() {
	ebegin "Stopping mwcollectd"
	start-stop-daemon --stop --quiet \
		--pidfile /var/run/mwcollect.pid \
		--exec /usr/sbin/mwcollectd
	eend $?
}

