# net-irc/unrealircd/files/unrealircd.rc
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-irc/unrealircd/files/unrealircd.rc,v 1.3 2004/07/14 23:48:05 agriffis Exp $

depend() {
	need net
	use dns
}

start() {
	ebegin "Starting unrealircd"
	start-stop-daemon --start --quiet --exec /usr/bin/unrealircd \
		--chuid ${UNREALIRCD_USER} -- ${UNREALIRCD_OPTS}
	eend $?
}

stop() {
	ebegin "Shutting down unrealircd"
	killall unrealircd
	eend $?
}
