# net-irc/unrealircd/files/unrealircd.rc
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-irc/unrealircd/files/unrealircd.rc,v 1.9 2011/12/21 18:25:03 binki Exp $

extra_started_commands="reload"

depend() {
	use dns net
	provide ircd
}

start() {
	ebegin "Starting unrealircd"
	start-stop-daemon --start --quiet --exec /usr/bin/unrealircd \
		--chuid ${UNREALIRCD_USER} -- ${UNREALIRCD_OPTS} &>/dev/null
	eend $?
}

stop() {
	ebegin "Shutting down unrealircd"
	start-stop-daemon --stop --quiet --exec /usr/bin/unrealircd
	eend $?
}

reload() {
	ebegin "Re-Loading unrealircd"
	killall -1 unrealircd
	eend $?
}
