# kde-base/kdenetwork/files/lisa
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/kde-base/kdenetwork/files/lisa,v 1.3 2004/03/04 18:45:39 vapier Exp $

depend() {
    use net
}

start() {
    ebegin "Starting lisa"
    start-stop-daemon --start -x _KDEDIR_/bin/lisa -q -- $LISACONF 1>/dev/null
    eend $? "Failed to start lisa"
}

stop() {
    ebegin "Stopping lisa"
    start-stop-daemon --stop -x _KDEDIR_/bin/lisa -q
    eend $? "Failed to stop lisa"
}
