# net-dns/pdnsd/files/pdnsd.online
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/net-dns/pdnsd/files/pdnsd.online,v 1.1 2004/02/02 16:59:39 dragonheart Exp $

depend() {
	need pdnsd net.${IFACE}
}

start() {
	ebegin "Starting pdnsd-online"
	/usr/sbin/pdnsd-ctl server all up
	eend $?
}

stop() {
	ebegin "Stopping pdnsd-online"
	/usr/sbin/pdnsd-ctl server all down
	eend $?
}

