# net-dns/pdnsd/files/pdnsd.online
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-dns/pdnsd/files/pdnsd.online,v 1.4 2004/07/14 23:26:19 agriffis Exp $

depend() {
	need pdnsd net.${IFACE}
}

start() {
	ebegin "Starting pdnsd-online"
	/usr/sbin/pdnsd-ctl server all up
	eend $?
}

stop() {
	ebegin "Stopping pdnsd-online"
	/usr/sbin/pdnsd-ctl server all down
	eend $?
}
