# gnustep-base/gnustep-base/files/gnustep
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/gnustep-base/gnustep-base/files/gnustep,v 1.1 2004/07/23 13:54:02 fafhrd Exp $

start() {
	ebegin "Starting GNUstep services"
	
	# Check for necessary files
	if [ -d /usr/GNUstep/System/Makefiles ]; then
		GNUSTEP_MAKEFILES="/usr/GNUstep/System/Makefiles"
        else
		exit 1
        fi

	# Source GNUstep makefile
	if [ -x $GNUSTEP_MAKEFILES/GNUstep.sh ]; then
		. $GNUSTEP_MAKEFILES/GNUstep.sh
	else
		exit 1
	fi

	# Start services
	if [ -x $GNUSTEP_SYSTEM_ROOT/Tools/$GNUSTEP_HOST_CPU/$GNUSTEP_HOST_OS/gdomap ]; then
		$GNUSTEP_SYSTEM_ROOT/Tools/$GNUSTEP_HOST_CPU/$GNUSTEP_HOST_OS/gdomap
	fi
	if [ -x $GNUSTEP_SYSTEM_ROOT/Tools/$GNUSTEP_HOST_CPU/$GNUSTEP_HOST_OS/$LIBRARY_COMBO/gdnc ]; then
		$GNUSTEP_SYSTEM_ROOT/Tools/$GNUSTEP_HOST_CPU/$GNUSTEP_HOST_OS/$LIBRARY_COMBO/gdnc
	fi
}

stop() {
	ebegin "Stopping GNUstep services"

	# Kill processes
	killall gdomap
	killall gdnc
}
