# sys-devel/distcc/files/2.11.1p/init
#!/sbin/runscript
# $Header: /var/cvsroot/gentoo-x86/sys-devel/distcc/files/2.11.1p/init,v 1.3 2004/03/21 21:44:24 lisa Exp $

depend() {
	need net
}

start() {
	rm -f ${DISTCCD_PIDFILE}
	ebegin "Starting distccd"
	TMPDIR="${TMPDIR}" \
	PATH="$(gcc-config --get-bin-path):${PATH}" \
	/sbin/start-stop-daemon --start --quiet --startas ${DISTCCD_EXEC} \
	--pidfile ${DISTCCD_PIDFILE} -- \
	--pid-file ${DISTCCD_PIDFILE} -N ${DISTCCD_NICE} --user distcc \
	${DISTCCD_OPTS}

	eend $?
}

stop() {
	ebegin "Stopping distccd"
	start-stop-daemon --stop --quiet --pidfile "${DISTCCD_PIDFILE}"
	rm -f "${DISTCCD_PIDFILE}"
	eend $?
}

