# app-admin/gkrellm/files/gkrellmd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/gkrellm/files/gkrellmd,v 1.4 2004/07/14 21:13:24 agriffis Exp $

depend() {
	need net
}

start() {
	ebegin "Starting GNU Krell Monitor daemon"
	start-stop-daemon --start --quiet --background --pidfile /var/run/gkrellmd.pid --make-pidfile --exec /usr/bin/gkrellmd
	eend $?
}

stop() {
	ebegin "Stopping GNU Krell Monitor daemon"
	start-stop-daemon --stop --quiet --pidfile /var/run/gkrellmd.pid --name gkrellmd
	eend $?
}
