# app-pda/plptools/files/psion
#!/sbin/runscript

#NB: Config is in /etc/conf.d/psion

depend() {
	use logger
}

start() {
    ebegin "Starting Psion Services"
    if test "$START_NCPD" = "yes" ; then
	ebegin "	ncpd"
	start-stop-daemon --start --quiet --exec /usr/sbin/ncpd \
			-- ${NCPD_ARGS} 1>&2
	eend
	if [ "$?" -eq 0 ]
	then
	    if test "$START_PLPNFSD" = "yes" ; then
		ebegin "	plpnfsd"
		start-stop-daemon --start --quiet --exec /usr/sbin/plpnfsd \
			-- ${PLPNFSD_ARGS} 1>&2
		eend
	    fi
	    if test "$START_PLPPRINTD" = "yes" ; then
		ebegin "	plpprintd"
		start-stop-daemon --start --quiet --exec /usr/sbin/plpprintd \
			-- ${PLPPRINTD_ARGS} 1>&2
		eend
	    fi
	fi
    fi
    eend $?
}

stop() {
    ebegin "Stopping Psion Services"
    if test "$START_PLPPRINTD" = "yes" ; then
	ebegin "	plpprintd"
	start-stop-daemon --start --quiet --exec /usr/sbin/plpprintd
	eend
    fi
    if test "$START_PLPNFSD" = "yes" ; then
	ebegin "	plpnfsd"
	echo stop > /mnt/psion/proc/exit
	eend
    fi
    if test "$START_NCPD" = "yes" ; then
	ebegin "	ncpd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/ncpd
	eend
    fi
    eend $?
}

