# sys-apps/athcool/files/athcool
#!/sbin/runscript
# $Header: /var/cvsroot/gentoo-x86/sys-apps/athcool/files/athcool,v 1.2 2003/08/10 13:42:09 lisa Exp $

start() {
	ebegin "Starting Athcool"
	/usr/sbin/athcool on &>/dev/null || return 1
	eend $?
}

stop() {
	ebegin "Stopping Athcool"
	/usr/sbin/athcool off &>/dev/null || return 1
	eend $?
}

