# sys-apps/athcool/files/athcool
#!/sbin/runscript
# $Header: /var/cvsroot/gentoo-x86/sys-apps/athcool/files/athcool,v 1.3 2004/07/18 04:31:56 dragonheart Exp $

start() {
	ebegin "Starting Athcool"
	/usr/sbin/athcool on &>/dev/null || return 1
	eend $?
}

stop() {
	ebegin "Stopping Athcool"
	/usr/sbin/athcool off &>/dev/null || return 1
	eend $?
}

