# dev-db/postgresql/files/7.2.3/postgresql
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# /space/gentoo/cvsroot/gentoo-x86/dev-db/mysql/files/mysql.rc6,v 1.1 2002/01/06 00:53:24 woodchip Exp

depend() {
	need net
}

start()	{
	ebegin "Starting postgres"
	start-stop-daemon  --start --quiet --background --chuid postgres --exec \
		/usr/bin/postmaster -- -D/var/lib/postgresql/data -N 1024 -B 2048
	eend $?
}

stop ()	{
	ebegin "Stopping postgres"
	start-stop-daemon --stop --quiet --pidfile=/var/lib/postgresql/data/postmaster.pid
	eend $?
}
