# net-dns/pdns/files/pdns
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/net-dns/pdns/files/pdns,v 1.1 2003/06/05 03:34:20 jhhudso Exp $

depend() {
	need net	
}

start() {
	ebegin "Starting PowerDNS"
	start-stop-daemon --start --quiet --exec /usr/bin/env -- /usr/sbin/pdns_server --daemon --guardian=yes
	eend $?
}

stop() {
	ebegin "Stopping PowerDNS"
	start-stop-daemon --stop --quiet --pidfile=/var/run/pdns.pid
	eend $?
}
