# app-emulation/vmware-workstation/files/3.0.0.1455-r1/vmware
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/app-emulation/vmware-workstation/files/3.0.0.1455-r1/vmware,v 1.2 2002/04/27 23:08:35 bangert Exp $


depend() {
	need localmount net
	use logger
}

start() {
	ebegin "Starting vmware services"
	test -x /etc/vmware/init.d/vmware || \
		eend 1 "vmware init script not found. Aborting" || return
	/etc/vmware/init.d/vmware start
	eend $? "Error Starting vmware services"
}

stop() {
	ebegin "Stopping vmware services"
	/etc/vmware/init.d/vmware stop
	eend $? "Error Stopping vmware services"
}
