# net-dns/dnsmasq/files/dnsmasq-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/net-dns/dnsmasq/files/dnsmasq-init,v 1.10 2007/03/23 10:44:02 uberlord Exp $

opts="reload"

depend() {
	provide dns
	need net
	use dbus
}

start() {
	ebegin "Starting dnsmasq"
	start-stop-daemon --start --exec /usr/sbin/dnsmasq \
	    --pidfile /var/run/dnsmasq.pid \
	    -- -x /var/run/dnsmasq.pid ${DNSMASQ_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping dnsmasq"
	start-stop-daemon --stop --exec /usr/sbin/dnsmasq \
	    --pidfile /var/run/dnsmasq.pid
	eend $?
}

reload() {
	ebegin "Reloading dnsmasq"
	start-stop-daemon --stop --oknodo --signal HUP \
	    --exec /usr/sbin/dnsmasq --pidfile /var/run/dnsmasq.pid
	eend $?
}
