# incoming/lcd/lcdproc
#!/sbin/runscript

port="ttyS0"
# which port is the LCD connected
driver="MtxOrb"
clients="C M T X"
# "C" (CPU usage),"G" (CPU Graph like Xload),"M" (Memory Usage),"X" (X-Load),
# "T" (Date, time, Uptime, OS version),"D" (Disk Stats),"B" (Battery Stats),"A" (About)
# for more Clients see "LCDproc Homepage"

start() {
	ebegin "Starting LCDProc"
	 ln -s /dev/$port /dev/lcd	
	 /usr/local/bin/LCDd -d $driver -b on -i off
	 /usr/local/bin/lcdproc $clients & 
	eend
}

stop() {
	ebegin "Stopping LCDProc"
	 killall lcdproc 
	 killall LCDd  
	 rm -f /dev/lcd
	eend
}
