# sys-cluster/heartbeat/files/heartbeat-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-cluster/heartbeat/files/heartbeat-init,v 1.5 2004/07/15 00:54:33 agriffis Exp $

HA_DIR=/etc/ha.d

depend() {
	use logger
	need net
}

# Run pre-startup script if it exists
RunStartStop() {
	[ -f $HA_DIR/resource.d/startstop ] && $HA_DIR/resource.d/startstop  "$@"
}

start() {
	ebegin "Starting Heartbeat"

	RunStartStop pre-start

	/usr/lib/heartbeat/heartbeat 2>&1
	ret=$?
	/usr/lib/heartbeat/ResourceManager verifyallidle

	RunStartStop post-start

	eend $ret
}

stop() {
	ebegin "Stopping Heartbeat"
	/usr/lib/heartbeat/heartbeat -k &> /dev/null
	eend $?
}

restart() {
	/usr/lib/heartbeat/heartbeat -r
}

status() {
	/usr/lib/heartbeat/heartbeat -s
}
