# net-misc/networkmanager/files/nm-interfaces
#!/sbin/runscript
# Copyright (c) 2012 Alexandre Rostovtsev
# Released under the 2-clause BSD license.

description="Virtual service for interfaces configured by NetworkManager. \
Do not start this service manually or add it to any runlevels!"

depend()
{
    need NetworkManager
    provide net
}

start()
{
    eerror "${RC_SVCNAME} should not be started manually or added to any runlevels."
    return 1
}

# Separate function to prevent warnings on starting service multiple times
mark_started()
{
    service_started && return 0
    ebegin "Marking NetworkManager network interfaces as started"
    mark_service_started
    eend 0
}

stop()
{
    [ "x$(LC_ALL=C nmcli -t -f state nm 2> /dev/null)" = "xconnected" ] && return 1
    ebegin "Marking NetworkManager network interfaces as stopped"
    eend 0
}
