# sys-power/powersoftplus/files/powersoftplus
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

start() {
	ebegin "Starting Ever UPS daemon"
	start-stop-daemon --start --quiet --exec /usr/bin/powersoftplus -- -q > /dev/null 2>&1
	eend $?
}

stop() {
	ebegin "Stopping Ever UPS daemon"
	start-stop-daemon --stop --retry 5 --pidfile /var/run/powersoftplus.pid
	eend $?
}
