# net-im/reaim/files/reaim
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-im/reaim/files/reaim,v 1.1 2004/06/15 19:06:30 squinky86 Exp $

depend() {
	need net iptables
}

start() {
	ebegin "Starting reaim"
	iptables -t nat -A PREROUTING -i eth0 -p tcp --dport 5190 -j REDIRECT --to-ports 5190
	start-stop-daemon --start --quiet --exec /usr/sbin/reaim
	eend $?
}

stop() {
	ebegin "Stopping reaim"
	iptables -t nat -D PREROUTING -i eth0 -p tcp --dport 5190 -j REDIRECT --to-ports 5190
	start-stop-daemon --stop --quiet --exec /usr/sbin/reaim
	eend $?
}
