# net-misc/asterisk/files/0.9.0/asterisk.rc6
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/asterisk/files/0.9.0/asterisk.rc6,v 1.3 2005/03/22 15:51:16 stkn Exp $

depend() {
	need net
	after zaptel
}

start() {
	ebegin "Starting asterisk PBX"
	/usr/sbin/asterisk ${OPTIONS}
	eend $?
}

stop() {
	ebegin "Stopping asterisk PBX"
	/usr/sbin/asterisk -r -x "STOP NOW" 2>/dev/null
	eend $?
}
