# net-misc/asterisk/files/0.9.0/asterisk.rc6
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/asterisk/files/0.9.0/asterisk.rc6,v 1.2 2004/07/14 23:54:29 agriffis Exp $

depend() {
	need net
}

start() {
	ebegin "Starting asterisk PBX"
	/usr/sbin/asterisk ${OPTIONS}
	eend $?
}

stop() {
	ebegin "Stopping asterisk PBX"
	/usr/sbin/asterisk -r -x "STOP NOW" 2>/dev/null
	eend $?
}
