# net-dialup/capi4k-utils/files/capi
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-dialup/capi4k-utils/files/capi,v 1.4 2004/07/14 22:39:53 agriffis Exp $

depend() {
	after isapnp
}

start() {
	if [ ! -e /etc/capi.conf ] ; then
		eerror "You're missing /etc/capi.conf (comes with a capi-driver)."
		eerror "Emerge net-dialup/fcpci if you are having an AVM Fritz!Card PCI"
		return 1
	fi
	/sbin/capiinit start
	eend 0
}

stop() {
	/sbin/capiinit stop
	eend $?
}
