# games-rpg/mangos/files/9999/mangos-realmd
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/games-rpg/mangos/files/9999/mangos-realmd,v 1.1 2008/10/04 07:38:26 trapni Exp $

depend() {
	need net
	use mysql
}

start() {
	ebegin "Starting MaNGOS Login Server (realmd)"
	start-stop-daemon --start --background -c mangos --exec /usr/share/mangos/bin/mangos-realmd \
		--make-pidfile --pidfile /var/run/mangos/realm.pid -- -c /etc/mangos/realmd.conf
	eend $?
}

stop() {
	ebegin "Stopping MaNGOS Login Server (realmd)"
	start-stop-daemon --stop -c mangos --exec /usr/share/mangos/bin/mangos-realmd \
		--pidfile /var/run/mangos/realm.pid
	eend $?
}
