# net-fs/idmapd/files/idmapd.confd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-fs/idmapd/files/idmapd.confd,v 1.2 2004/07/14 23:34:24 agriffis Exp $

# Options to pass to idmapd
#IDMAPD_OPTS="-d domain"
