# app-admin/fam-oss/files/fam
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/app-admin/fam-oss/files/fam,v 1.2 2001/12/06 20:17:01 azarah Exp $

depend() {
	need portmap
}


start() {
	ebegin "Starting fam"
	#for some reason dont want to daemonize without the '&'
	start-stop-daemon --start --quiet --exec /usr/bin/fam \
		-- -T 0 -c /etc/fam.conf &
	eend $?
}

stop() {
	ebegin "Stopping fam"
	start-stop-daemon --stop --quiet --exec /usr/bin/fam
	eend $?
}
