# app-admin/fam-oss/files/fam
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/fam-oss/files/fam,v 1.4 2002/11/09 10:02:23 vapier Exp $

depend() {
	need portmap
}

start() {
	ebegin "Starting fam"
	start-stop-daemon --start --quiet --exec /usr/bin/fam --background \
		-- -T 0 -c /etc/fam.conf
	eend $?
}

stop() {
	ebegin "Stopping fam"
	start-stop-daemon --stop --quiet --exec /usr/bin/fam
	eend $?
}
