# sys-apps/halevt/files/halevt
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later

depend() {
	need hald
}

start() {
	ebegin "Starting Automounter"
	start-stop-daemon --start --pidfile /var/lib/run/halevt/halevt.pid --exec /usr/bin/halevt -- -u root -g plugdev > /dev/null
	eend $?
}

stop() {
	ebegin "Stopping Automounter"
	start-stop-daemon --stop --pidfile /var/lib/run/halevt/halevt.pid
	eend $?
}
