# net-firewall/firehol/files/firehol.initrd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-firewall/firehol/files/firehol.initrd,v 1.2 2004/07/14 23:28:26 agriffis Exp $

opts="start stop"

depend() {
	need net
	provide firewall
}

start() {
	ebegin "Starting FireHOL"
	/usr/sbin/firehol start 1>/dev/null
	eend $? 
}

stop() {
	ebegin "Stopping FireHOL"
	/usr/sbin/firehol stop 1>/dev/null
	eend $?
}

