# gnustep-base/gnustep-env/files/gnustep.runscript-0.1.4
#!/sbin/runscript
#
# (c) 2003 Dennis Leeuw
# License: GPL version 2 or any newer version
#
# Modifications, truncation, and gentoo-ification:
# Copyright 2004 Gentoo Technologies, Inc.

GNUSTEP_SYSTEM_ROOT="/usr/GNUstep/System"
. ${GNUSTEP_SYSTEM_ROOT}/Library/Makefiles/GNUstep.sh
TOOLS="${GNUSTEP_SYSTEM_ROOT}/Tools"

start() {
	ebegin "Starting gnustep daemons..."

	if [ -x $TOOLS/gdomap ]; then 
		ebegin "\tgdomap "
		GNUSTEP_USER_ROOT=/var/run/GNUstep \
		GNUSTEP_DEFAULTS_ROOT=/var/run/GNUstep \
		$TOOLS/gdomap
	fi
	if [ -x $TOOLS/gdnc ]; then
		ebegin "\tgdnc "
		GNUSTEP_USER_ROOT=/var/run/GNUstep \
		GNUSTEP_DEFAULTS_ROOT=/var/run/GNUstep \
		$TOOLS/gdnc
	fi
	if [ -x $TOOLS/gnustep_sndd ]; then 
		ebegin "\tgnustep_sndd"
		GNUSTEP_USER_ROOT=/var/run/GNUstep \
		GNUSTEP_DEFAULTS_ROOT=/var/run/GNUstep \
		$TOOLS/gnustep_sndd
	fi

	eend $?
}

stop() {
	ebegin "Stopping gnustep daemons..."

	ebegin "\tgdomap "
	killall gdomap
	
	ebegin "\tgdnc "
	killall gdnc

	ebegin "\tgnustep_sndd"
	killall gnustep_sndd

	eend $?
}

