# x11-misc/keytouch/files/keytouch
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/x11-misc/keytouch/files/keytouch,v 1.1 2006/10/06 13:19:11 nyhm Exp $

start() {
	ebegin "Starting keyTouch"
	/usr/bin/keytouch-init
	eend
}
