# net-nds/ypserv/files/ypserv-initd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-nds/ypserv/files/ypserv-initd,v 1.2 2004/07/15 00:19:17 agriffis Exp $

depend() {
	need net
}

start() {
	ebegin "Starting NIS Server"
	start-stop-daemon --start --quiet --exec /usr/sbin/ypserv
	result=${?}

	if [ ${PW_SERVER} = "yes" ]; then
		start-stop-daemon --start --quiet --exec /usr/sbin/rpc.yppasswdd ${PW_SERVER_OPTIONS}
		result=$(( $result + $? ))
	fi

	if [ ${XFR_SERVER} = "yes" ]; then
		start-stop-daemon --start --quiet --exec /usr/sbin/rpc.ypxfrd
		result=$(( $result + $? ))
	fi

	eend ${result}
}

stop() {
	ebegin "Stopping NIS Server"
	start-stop-daemon --stop --quiet --exec /usr/sbin/ypserv
	result=${?}

	if [ ${PW_SERVER} = "yes" ]; then
		start-stop-daemon --stop --quiet --exec /usr/sbin/rpc.yppasswdd
		result=$(( $result + $? ))
	fi

	if [ ${XFR_SERVER} = "yes" ]; then
		start-stop-daemon --stop --quiet --exec /usr/sbin/rpc.ypxfrd
		result=$(( $result + $? ))
	fi

	eend ${result}
}
