# net-misc/dante/files/dante-sockd-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/dante/files/dante-sockd-init,v 1.4 2004/07/14 23:56:28 agriffis Exp $

depend() {
	need net
}

checkconfig() {
	if [ ! -f /etc/socks/sockd.conf ] ; then
		eerror "You need to setup /etc/socks/sockd.conf first"
		eerror "Examples are in /usr/share/doc/dante/example"
		return 1
	fi
}

start() {
	checkconfig || return 1
	ebegin "Starting dante sockd"
	start-stop-daemon --start --quiet --pidfile /var/run/sockd.pid \
		--make-pidfile --background --exec /usr/sbin/sockd
	eend $? "Failed to start sockd"
}

stop() {
	ebegin "Stopping dante sockd"
	start-stop-daemon --stop --quiet --pidfile /var/run/sockd.pid
	eend $? "Failed to stop sockd"
	# clean stale pidfile
	[ -f /var/run/sockd.pid ] && rm -f /var/run/sockd.pid
}
