# x11-misc/mkxf86config/files/mkxf86config-init
#!/sbin/runscript
 
depend() {
    before local
}
 
start() {
    ebegin "Creating X Configuration"
    /usr/sbin/mkxf86config.sh
	eend $?
}
