# x11-misc/mkxf86config/files/mkxf86config-init
#!/sbin/runscript
 
depend() {
    need
}
 
start() {
    ebegin "Create XF86Config"
    /usr/sbin/mkxf86config.sh
}
