# net-misc/powerd/files/powerd-init
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later

depend() {
	after local
}

start() {
	ebegin "Starting powerd"
	start-stop-daemon --start --exec /usr/sbin/powerd
	eend $?
}

stop() {
	ebegin "Stopping powerd"
	killall powerd
	eend $?
}

