# media-video/flumotion/files/flumotion-init
#!/sbin/runscript

flumotion_manager=/usr/bin/flumotion-manager
flumotion_worker=/usr/bin/flumotion-worker
flumotion_manager_pidfile=/var/run/flumotion/manager.default.pid
flumotion_worker_pidfile=/var/run/flumotion/worker.localhost.pid

depend() {
    need net
}

start() {
    ebegin "Starting Flumotion Manager"
    HOME=/usr/share/flumotion start-stop-daemon --start --chuid flumotion:flumotion --pidfile ${flumotion_manager_pidfile} --exec ${flumotion_manager} -- -D /etc/flumotion/managers/default/planet.xml 
    eend $?
    ebegin "Starting Flumotion Worker"
    HOME=/usr/share/flumotion start-stop-daemon --start --chuid flumotion:flumotion --pidfile ${flumotion_worker_pidfile} --exec ${flumotion_worker} -- -D /etc/flumotion/workers/default.xml
    eend $?
}

stop() {
    ebegin "Stopping Flumotion Worker"
    start-stop-daemon --stop --pidfile ${flumotion_worker_pidfile}
    eend $?
    ebegin "Stopping Flumotion Manager"
    start-stop-daemon --stop --pidfile ${flumotion_manager_pidfile}
    eend $?
}
