# app-admin/verynice/files/verynice
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/verynice/files/verynice,v 1.2 2004/07/14 21:21:11 agriffis Exp $

start() {
	ebegin "Starting verynice"
	/usr/sbin/verynice -d /var/run/verynice.pid
	eend $?
}

stop() {
	ebegin "Stopping verynice"
	start-stop-daemon --stop --quiet --pidfile=/var/run/verynice.pid
	eend $?
}
