# sys-fs/multipath-tools/files/rc-multipathd
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-fs/multipath-tools/files/rc-multipathd,v 1.2 2005/02/21 15:10:57 azarah Exp $

depend() {
	need localmount
}

start() {
	ebegin "Starting multipathd"
	start-stop-daemon --start --quiet --exec /usr/sbin/multipathd
	eend $?
}

stop() {
	ebegin "Stopping multipathd"
	start-stop-daemon --stop --quiet --pid /var/run/multipathd.pid
	eend $?
}
