# app-misc/welcome2l/files/welcome2l.initscript
#!/sbin/runscript

depend() {
	need local
}

start() {
	ebegin "Starting Welcome2L"
	/usr/bin/Welcome2L -getty -msg"Gentoo Linux" > /etc/issue
	eend $result
}
