# net-wireless/irda-utils/files/irda-usb.rules
KERNEL=="irda[0-9]*", RUN+="/lib/udev/irda-usb.sh"
