# net-dns/dns2go/files/dns2go.rc6
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-dns/dns2go/files/dns2go.rc6,v 1.4 2004/07/14 23:22:19 agriffis Exp $

depend() {
	need net
}

start() {
	ebegin "Starting DNS2GO"
	/usr/bin/dns2go
	eend $?
}

stop() {
	ebegin "Stopping DNS2GO"
	/usr/bin/dns2go -k
	eend $?
}
