# net-analyzer/vnstat/files/vnstat.cron
#!/bin/bash

# uncomment this line to enable this cron script
exit 0

if [ `ls /var/lib/vnstat/ | wc -l` -ge 1 ] ; then
	exec /usr/bin/vnstat -u
fi
