# net-analyzer/net-snmp/files/snmpd.rc6
#!/sbin/runscript

depend() {
	need net
	use logger
}

checkconfig() {
	if [ ! -e /etc/snmp/snmpd.conf ] ; then
		eerror "You need an /etc/snmp/snmpd.conf to run snmpd"
		return 1
	fi
}

start() {
	checkconfig || return 1
	ebegin "Starting net-snmpd"
	start-stop-daemon --start --quiet --exec /usr/sbin/snmpd \
		-- -C -c /etc/snmp/snmpd.conf ${SNMPD_FLAGS}
	eend $?
}

stop() {
	ebegin "Stopping net-snmpd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/snmpd
	eend $?
}
