# app-misc/lcdproc/files/LCDd
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/app-misc/lcdproc/files/LCDd,v 1.1 2003/01/07 16:59:04 aliz Exp $

depend() {
	need net
}

start() {
	ebegin "Starting LCDd"
	/usr/sbin/LCDd -c /etc/LCDd.conf
	eend $?
}

stop() {
	ebegin "Stopping LCDd"
	killall LCDd
	eend $?
}
