# app-misc/lcdproc/files/LCDd
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/app-misc/lcdproc/files/LCDd,v 1.4 2003/02/14 22:36:43 vapier Exp $

depend() {
	need net
}

start() {
	ebegin "Starting LCDd"
	start-stop-daemon --start --quiet -m --pidfile /var/run/LCDd.pid \
		--background --exec /usr/sbin/LCDd -- -s -f -c /etc/LCDd.conf
	sleep 1
	eend $?
}

stop() {
	ebegin "Stopping LCDd"
	start-stop-daemon --stop --quiet --pidfile /var/run/LCDd.pid
	eend $?
}
