# app-misc/lcdproc/files/lcdproc
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/app-misc/lcdproc/files/lcdproc,v 1.5 2003/01/07 16:58:28 aliz Exp $

start() {
	ebegin "Starting LCDProc"
	source /etc/lcdproc.conf
	ebegin "Starting lcdproc"
	/usr/bin/lcdproc -d ${SCREENS}
	eend $?
}

stop() {
	ebegin "Stopping LCDProc"
	killall lcdproc 
	eend $?
}
