# app-misc/lcdproc/files/lcdproc
#!/sbin/runscript

# Please don't change these variables.  They are listed here as
# defaults for reference.  Make your modifications in
# /etc/conf.d/lcdproc instead.
port="ttyS0"
driver="MtxOrb"
clients="C M T X"

start() {
	ebegin "Starting LCDProc"
	 ln -s /dev/$port /dev/lcd	
	 /usr/local/bin/LCDd -d $driver -b on -i off
	 /usr/local/bin/lcdproc $clients & 
	eend
}

stop() {
	ebegin "Stopping LCDProc"
	 killall lcdproc 
	 killall LCDd  
	 rm -f /dev/lcd
	eend
}
