# app-misc/lcdproc/files/lcdproc
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-misc/lcdproc/files/lcdproc,v 1.9 2004/02/13 19:19:04 vapier Exp $

depend(){
	need LCDd 
}
start() {
	source /etc/lcdproc.conf
	ebegin "Starting lcdproc"
	start-stop-daemon --start --background \
		--exec /usr/bin/lcdproc -- ${SCREENS}
	eend $?
}

stop() {
	ebegin "Stopping lcdproc"
	start-stop-daemon --stop --exec /usr/bin/lcdproc
	eend $?
}
