# net-mail/qmail/files/1.03-r13/qmail-control
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-mail/qmail/files/1.03-r13/qmail-control,v 1.3 2004/03/06 03:20:50 vapier Exp $

opts="start stop restart reload"

start() {
    ebegin "Starting qmail mta ..."
    svc -u /service/qmail-{smtpd,send,qmqpd,qmtpd}
    eend $?
    ebegin "Starting qmail mta logging ..."
    svc -u /service/qmail-{smtpd,send,qmqpd,qmtpd}/log
    eend $?
}

stop() {
    ebegin "Stopping qmail mta ..."
    svc -d /service/qmail-{smtpd,send,qmqpd,qmtpd}
    eend $?
    ebegin "Stopping qmail mta logging ..."
    svc -d /service/qmail-{smtpd,send,qmqpd,qmtpd}/log
    eend $?
}

reload() {
    ebegin "Reloading 'locals' and 'virtualdomains' control files."
    svc -h /service/qmail-send
    eend $?
}
