# app-admin/gps/files/rgpsp
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/gps/files/rgpsp,v 1.3 2004/07/14 21:13:45 agriffis Exp $

depend() {
	need net
}

start() {
	ebegin "Starting Remote Graphical Process Statistics daemon"
	start-stop-daemon --start --quiet --exec /usr/bin/rgpsp
	eend $?
}

stop() {
	ebegin "Stopping Remote Graphical Process Statistics daemon"
	start-stop-daemon --stop --quiet --name rgpsp
	eend $?
}
