# sys-kernel/kerneloops/files/kerneloops
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-kernel/kerneloops/files/kerneloops,v 1.2 2008/01/08 23:45:41 gregkh Exp $

depend() {
	need net
}

exec="/usr/sbin/kerneloops"
prog=$(basename $exec)
conf="/etc/kerneloops.conf"
pidfile="/var/run/${prog}.pid"

checkconfig() {
	if [ ! -e $conf ] ; then
		eerror "You need an ${conf} file to run $prog"
		return 1
	fi
}

start() {
	checkconfig || return 1
	ebegin "Starting Kerneloops Monitor daemon"
	start-stop-daemon --start --quiet --background --pidfile ${pidfile} --make-pidfile --exec ${exec} -- --nodaemon
	eend $?
}

stop() {
	ebegin "Stopping Kerneloops Monitor daemon"
	start-stop-daemon --stop --quiet --pidfile ${pidfile} --name ${prog}
	eend $?
}
