# net-analyzer/zabbix-agent/files/1.1_alpha7/init.d/zabbix-agentd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-analyzer/zabbix-agent/files/1.1_alpha7/init.d/zabbix-agentd,v 1.2 2005/05/10 10:46:11 swegener Exp $

opts="${opts} reload"

depend() {
		need net
		provide zabbix-agent
}

start() {
	ebegin "Starting Zabbix agent"
	start-stop-daemon --start --exec /usr/sbin/zabbix_agentd
	eend
}
stop() {
	ebegin "Stopping Zabbix agent"
	start-stop-daemon  --stop --pid $PID
	eend
}
