# sys-power/nut/files/upsdrv.rc6-r1
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-power/nut/files/upsdrv.rc6-r1,v 1.1 2006/05/15 23:43:08 robbat2 Exp $

depend() {
	need net
	before upsd
}

start() {
	ebegin "Starting UPS drivers"
	/usr/sbin/upsdrvctl start
	eend $? "Failed to start UPS drivers!"
}

stop() {
	ebegin "Stopping UPS drivers"
	/usr/sbin/upsdrvctl stop
	eend $? "Failed to stop UPS drivers!"
}
