# dev-db/xindice/files/xindice
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-db/xindice/files/xindice,v 1.3 2004/07/14 21:42:52 agriffis Exp $

depend() {
	need net
}

start() {
	ebegin "Starting xindice..."
	start-stop-daemon --pidfile /var/run/xindice/xindice.pid\
		--start --background --user xindice --chuid xindice\
		--exec /usr/lib/xindice/start	
	eend $?
}


stop() {
	ebegin "Stopping xindice..."
	start-stop-daemon --pidfile /var/run/xindice/xindice.pid --stop	--retry 10
	eend $?
}
