# net-www/wwwoffle/files/wwwoffled
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-www/wwwoffle/files/wwwoffled,v 1.7 2004/03/06 04:03:00 vapier Exp $

EXE=`which wwwoffled`

depend() {
	need localmount
}

start() {
	ebegin "Starting wwwoffled"
	start-stop-daemon --start --quiet --exec $EXE \
		-- -c /etc/wwwoffle.conf
	eend $?
}

stop() {
	ebegin "Stopping wwwoffled"
	start-stop-daemon --stop --quiet --exec $EXE
	eend $?
}
