# net-www/wwwoffle/files/wwwoffled
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-www/wwwoffle/files/wwwoffled,v 1.8 2004/04/25 00:44:10 dragonheart Exp $

depend() {
	need localmount
}

start() {
	ebegin "Starting wwwoffled"
	start-stop-daemon --start --quiet --exec `which wwwoffled` \
		-- -c /etc/wwwoffle.conf
	eend $?
}

stop() {
	ebegin "Stopping wwwoffled"
	start-stop-daemon --stop --quiet --name wwwoffled
	eend $?
}
