# net-www/wwwoffle/files/wwwoffled-online
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/net-www/wwwoffle/files/wwwoffled-online,v 1.2 2001/12/06 20:17:02 azarah Exp $

depend() {
    need wwwoffled net
}

start() {
    ebegin "Starting wwwoffled-online"
    wwwoffle -online
    eend $? "Started wwwoffled-online." "Error starting wwwoffled-online."
}

stop() {
    ebegin "Stopping wwwoffled-online"
    wwwoffle -offline
    eend $? "Stopped wwwoffled-online." "Error stopping wwwoffled-online."
}

