# net-www/wwwoffle/files/wwwoffled-online
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/net-www/wwwoffle/files/wwwoffled-online,v 1.6 2003/08/03 15:09:45 wmertens Exp $

depend() {
    need wwwoffled net
}

start() {
    ebegin "Starting wwwoffled-online"
    wwwoffle -online -c /etc/wwwoffle.conf
    wwwoffle -fetch -c /etc/wwwoffle.conf
    eend $?
}

stop() {
    ebegin "Stopping wwwoffled-online"
    wwwoffle -offline -c /etc/wwwoffle.conf
    eend $?
}

