# net-wireless/wimax/files/wimax
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-wireless/wimax/files/wimax,v 1.2 2010/07/26 22:13:27 alexxy Exp $

description="WiMAX Network Service"

depend() {
	need localmount
	after bootmisc
}

start() {
	ebegin "Starting WiMAX Network Service"
	[ -e /var/run/wimaxd.pid ] && rm -f /var/run/wimaxd.pid
	start-stop-daemon --start --quiet --pidfile /var/run/wimax.pid --exec /usr/bin/wimaxd -b
	eend $?
}

stop() {
	ebegin "Stoping WiMAX Network Service"
	start-stop-daemon --stop --quiet --pidfile /var/run/wimax.pid
	eend $?
}
