# net-mail/qmail-ldap/files/qmail-control
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/net-mail/qmail-ldap/files/qmail-control,v 1.1 2003/07/26 02:05:06 raker Exp $

opts="start stop restart reload"

start() {
    ebegin "Starting qmail mta ..."
    svc -u /var/qmail/supervise/*
    eend $?
    ebegin "Starting qmail mta logging ..."
    svc -u /var/qmail/supervise/*/log
    eend $?
}

stop() {
    ebegin "Stopping qmail mta ..."
    svc -d /var/qmail/supervise/*
    eend $?
    ebegin "Stopping qmail mta logging ..."
    svc -d /var/qmail/supervise/*/log
    eend $?
}

reload() {
    ebegin "Reloading 'locals' and 'virtualdomains' control files."
    svc -h /var/qmail/supervise/qmail-send
    eend $?
}
