# sys-fs/lufs/files/lufs-0.7.0-init
#!/sbin/runscript

# lufsd
#
# description: LUFS Userspace Daemon
#
# processname: lufs
# chkconfig: - 91 91

depend () {
        need localmount
}

start () {
        echo -n "Starting lufsd... "
        su - lufs -c "lufsd -d 2>~/lufsd.err >~/lufsd.log"
        echo
}

stop () {
        echo -n "Shutting down lufsd... "
        killall lufsd
        echo
}
