# sys-fs/lufs/files/lufs-0.9.3-init
#!/sbin/runscript

# lufsd
#
# description: LUFS Userspace Daemon
#
# processname: lufs
# chkconfig: - 91 91
#
# lufsd is started by the lufs user, and logs are kept in its homedir.
#

depend () {
        need localmount
}

start () {
        ebegin "Starting lufsd"
        su - lufs -c "lufsd -d 2>~/lufsd.err >~/lufsd.log"
        eend $? "Error starting lufsd"
}

stop () {
        ebegin "Stopping lufsd"
        pkill lufsd
        eend $? "Error stopping lufsd"
}
