# net-mail/teapop/files/teapop-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-mail/teapop/files/teapop-init,v 1.2 2004/03/06 03:24:04 vapier Exp $

depend() {
	need net
	use logger
	after *
}

checkconfig() {
	if [ ! -e /etc/teapop.passwd ] ; then
		eerror "You need an /etc/teapop.passwd file to run teapop"
		return 1
	fi
}

start() {
	checkconfig || return 1
	ebegin "Starting teapop"
	start-stop-daemon --start --quiet \
		--pidfile /var/run/teapop.pid --exec /usr/sbin/teapop -- ${TEAPOP_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping teapop"
	start-stop-daemon --stop --quiet \
		--pidfile /var/run/teapop.pid --exec /usr/sbin/teapop
	eend $?
}

restart() {
	stop
	start
}
