# sys-apps/autospeedstep/files/autospeedstep
#!/sbin/runscript

start() {
	ebegin "Starting autospeedstep"
	/usr/sbin/autospeedstep
	eend $? "Failed to start autospeedstep"
}
