# net-misc/zebedee/files/zebedee
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/zebedee/files/zebedee,v 1.2 2004/03/06 03:40:48 vapier Exp $

depend() {
	need net
	use dns logger
}

start() {
	ebegin "Starting Zebedee"
	start-stop-daemon --quiet --start --startas /usr/bin/zebedee \
		-m -b -p /var/run/zebedee.pid \
		-- -f /etc/zebedee/server.zbd
	eend $? "Failed to start Zebedee"
}

stop() {
	ebegin "Stopping Zebedee"
	start-stop-daemon --quiet --stop \
	-p /var/run/zebedee.pid
	eend $? "Failed to stop Zebedee"

	# clean stale pidfile
	[ -f /var/run/zebedee.pid ] && rm -f /var/run/zebedee.pid
}

restart() {
	ebegin "Restarting Zebedee"
	svc_stop
	svc_start
	eend $? "Failed to restart Zebedee"
}
