# net-misc/linux-identd/files/identd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/linux-identd/files/identd,v 1.3 2004/03/06 03:36:03 vapier Exp $

depend() {
	need net
}

start() {
	ebegin "Starting identd"
	start-stop-daemon --start --quiet --name identd \
		--startas /usr/sbin/identd -- -d
	eend $?
}

stop() {
	ebegin "Stopping identd"
	start-stop-daemon --stop --quiet --name identd
	eend $?
}
