# sys-devel/distcc/files/2.8/init
#!/sbin/runscript
# $Header: /var/cvsroot/gentoo-x86/sys-devel/distcc/files/2.8/init,v 1.2 2003/08/08 20:31:11 lisa Exp $

depend() {
	need net
}

start() {
	rm -f ${DISTCCD_PIDFILE}
	ebegin "Starting distccd"
	DISTCC_TMPDIR="${DISTCC_TMPDIR}" \
	PATH="$(gcc-config --get-bin-path):${PATH}" \
	    /sbin/start-stop-daemon --start --quiet \
	        --startas ${DISTCCD_EXEC} \
		-c distcc --pidfile ${DISTCCD_PIDFILE} \
		-- --pid-file ${DISTCCD_PIDFILE} ${DISTCCD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping distccd"
	start-stop-daemon --stop --quiet --pidfile "${DISTCCD_PIDFILE}"
	rm -f "${DISTCCD_PIDFILE}"
	eend $?
}

