# sys-apps/vcron/files/vcron.rc6
#!/sbin/runscript

need() {
	clock localmount
}

start() {
	ebegin "Starting vcron"
	start-stop-daemon --start --quiet --exec /usr/bin/cron
	eend $?
}

stop() {
	ebegin "Stopping vcron"
	start-stop-daemon --stop --quiet --exec /usr/bin/cron
	eend $?
}
