# sys-apps/vcron/files/vcron.rc6
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/sys-apps/vcron/files/vcron.rc6,v 1.2 2001/12/06 20:17:03 azarah Exp $

need() {
	clock localmount
}

start() {
	ebegin "Starting vcron"
	start-stop-daemon --start --quiet --exec /usr/bin/cron
	eend $?
}

stop() {
	ebegin "Stopping vcron"
	start-stop-daemon --stop --quiet --exec /usr/bin/cron
	eend $?
}
