# net-irc/atheme-services/files/atheme-services.initd.in
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-irc/atheme-services/files/atheme-services.initd.in,v 1.1 2011/06/09 04:38:07 binki Exp $

opts="reload"

depend() {
	need net
	use ircd
	provide irc-services
}

start() {
	ATHEME_SERVICES_RUNDIR="@GENTOO_PORTAGE_EPREFIX@"/var/run/atheme
	if ! [ -d "${ATHEME_SERVICES_RUNDIR}" ]; then
	    ebegin "Creating \`${ATHEME_SERVICES_RUNDIR}'"

	    ATHEME_INSTALL_O="-o atheme-services"
	    [ "@GENTOO_PORTAGE_EPREFIX@" ] && ATHEME_INSTALL_O=

	    install ${ATHEME_INSTALL_O} -d "${ATHEME_SERVICES_RUNDIR}"
	    eend $?
	fi

	ebegin "Starting ${SVCNAME}"

	ATHEME_SSD_CHUID="--chuid atheme-services"
	[ "@GENTOO_PORTAGE_EPREFIX@" ] && ATHEME_SSD_CHUID=

	start-stop-daemon --start ${ATHEME_SSD_CHUID} --exec "@GENTOO_PORTAGE_EPREFIX@"/usr/bin/atheme-services
	eend $?
}

reload() {
	ebegin "Reloading ${SVCNAME}"
	start-stop-daemon --signal HUP --pidfile "@GENTOO_PORTAGE_EPREFIX@"/var/run/atheme/atheme.pid
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --quiet --pidfile "@GENTOO_PORTAGE_EPREFIX@"/var/run/atheme/atheme.pid
	eend $?
}
