# mail-filter/dkim-milter/files/dkim-filter-2.4.1-init
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/mail-filter/dkim-milter/files/dkim-filter-2.4.1-init,v 1.1 2007/12/28 10:09:57 dragonheart Exp $

CONFFILE=/etc/mail/dkim-filter/${SVCNAME}.conf
PIDFILE=/var/run/dkim-filter/${SVCNAME}.pid

depend() {
	use dns logger net
}

check_cfg() {

	if [ ! -f "${CONFFILE}" ] ; then
		eerror "Configuration file ${CONFFILE} is missing"
		return 1
	fi

	if egrep -q '^[[:space:]]*Background[[:space:]]+no' "${CONFFILE}" ; then
		eerror "${SVCNAME} service cannot run with Background key set to yes!"
		return 1
	fi

	if egrep -q '^[[:space:]]*PidFile[[:space:]]' "${CONFFILE}"; then
		eerror "Please remove PidFile key from your configuration file!"
		eerror "${SVCNAME} is automatically configured to use ${PIDFILE} as pid file."
		return 1
	fi
}

start() {
	check_cfg || return 1

	ebegin "Starting DKIM filter"
	start-stop-daemon --start --pidfile "${PIDFILE}" \
		--exec /usr/bin/dkim-filter -- -x "${CONFFILE}" -P "${PIDFILE}"
	eend $?
}

stop() {
	ebegin "Stopping DKIM filter"
	start-stop-daemon --stop --pidfile "${PIDFILE}"
	eend $?
}
