# net-misc/vpnc/files/vpnc
#!/sbin/runscript

depend() {
	need net
}

start() {
  ebegin "Starting Virtual Private Network"
  start-stop-daemon --start --exec /usr/bin/vpnc \
    --pidfile /var/run/vpnc/pid
  eend $?
}

stop() {
  ebegin "Stopping Virtual Private Network"
  start-stop-daemon --stop --exec /usr/bin/vpnc-disconnect \
    --pidfile /var/run/vpnc/pid
  eend $?
}
