# x11-base/xfree/files/4.1.0-r6/xfs.start
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/x11-base/xfree/files/4.1.0-r6/xfs.start,v 1.2 2003/02/14 23:29:31 vapier Exp $

depend() {
	use logger
}

start() {
	ebegin "Starting X Font Server"
	if [ "`grep xfs /etc/passwd`" ] ; then
		start-stop-daemon --start --quiet --exec /usr/X11R6/bin/xfs \
			-- -port -1 -daemon -config /etc/X11/fs/config \
				-droppriv -user xfs 1>&2
	else
		start-stop-daemon --start --quiet --exec /usr/X11R6/bin/xfs \
			-- -port -1 -daemon -config /etc/X11/fs/config 1>&2
	fi
	eend $?
}

stop() {
	ebegin "Stopping X Font Server"
	start-stop-daemon --stop --quiet --exec /usr/X11R6/bin/xfs 1>&2
	rm -rf /tmp/.font-unix
	eend $?
}
