# sys-devel/distcc/files/2.7.1/init
#!/sbin/runscript

depend() {
	need net
}

start() {
	rm -f ${DISTCCD_PIDFILE}
	ebegin "Starting distccd"
	PATH="$(gcc-config --get-bin-path):${PATH}" \
	    /sbin/start-stop-daemon --start --quiet \
	        --startas ${DISTCCD_EXEC} \
		-c distcc --pidfile ${DISTCCD_PIDFILE} \
		-- --pid-file ${DISTCCD_PIDFILE} ${DISTCCD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping distccd"
	start-stop-daemon --stop --quiet --pidfile "${DISTCCD_PIDFILE}"
	rm -f "${DISTCCD_PIDFILE}"
	eend $?
}

