# app-emulation/xen-tools/files/xendomains-init
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-emulation/xen-tools/files/xendomains-init,v 1.1 2006/03/24 23:55:01 agriffis Exp $

opts="start stop status restart"

depend() {
	need xend
	after dhcp
}

get_domname() {
	sed -ne 's/^[ \t]*name[ \t]*=[ \t]*"\([^"]*\)"/\1/p' $1
}

is_running() {
	[ -n "`/usr/sbin/xm list | grep "^${1} "`" ]
}

start() {
	einfo "Starting ${AUTODIR} Xen domains"
	if [[ ${SCREEN} == "yes" ]]; then
		screen -d -m -S xen -t dom0
		screen -r xen -X zombie dr
		logrotate -f /usr/share/xen/xen-consoles-logrotate
		screen -r xen -X logfile /var/log/xen-consoles/%t
		screen -r xen -X logfile flush 1
		screen -r xen -X deflog on
	fi
	# Create all domains with config files in AUTODIR.
	for dom in $(ls ${AUTODIR}/* 2>/dev/null); do
		name=$(get_domname ${dom})
		if ! is_running ${name} ; then
			ebegin "  Starting domain ${name}"
			if [[ ${SCREEN} == "yes" ]]; then
				screen -r xen -X screen -t ${name} xm create ${dom} -c
			else
				xm create --quiet ${dom}
			fi
			eend $?
		else
			einfo "  Not Starting domain ${name} - allready running"
		fi
	done
}

stop() {
	einfo "Shutting down ${AUTODIR} Xen domains"
	# Stop all domains with config files in AUTODIR.
	for dom in $(ls ${AUTODIR}/* 2>/dev/null); do
		name=$(get_domname ${dom})
		if is_running ${name} ; then
			ebegin "  Stopping domain ${name}"
				xm shutdown --wait ${name} >/dev/null
			eend $?
		else
			einfo "  Not Stopping domain ${name} - not running"
		fi
	done
	if [[ ${SCREEN} == "yes" ]]; then
		screen -r xen -X quit
	fi
}

status() {
	/usr/sbin/xm list
}
