# dev-libs/cyrus-sasl/files/pwcheck.rc6
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/dev-libs/cyrus-sasl/files/pwcheck.rc6,v 1.2 2003/02/14 22:40:41 vapier Exp $

depend() {
	need localmount
	use logger
}

start() {
	ebegin "Starting sasl pwcheck daemon"
	start-stop-daemon --start --quiet --oknodo --exec /usr/sbin/pwcheck
	eend $?
}

stop() {
	ebegin "Stopping sasl pwcheck daemon"
	start-stop-daemon --stop --quiet --oknodo --exec /usr/sbin/pwcheck
	eend $?
}
