# dev-libs/cyrus-sasl/files/saslauthd2.rc6
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/dev-libs/cyrus-sasl/files/saslauthd2.rc6,v 1.4 2004/01/14 19:42:47 max Exp $

depend() {
	need net
}

start() {
	ebegin "Starting saslauthd"
	start-stop-daemon --start --quiet --exec /usr/sbin/saslauthd \
		-- ${SASLAUTHD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping saslauthd"
	start-stop-daemon --stop --quiet --pidfile /var/lib/sasl2/saslauthd.pid
	eend $?
}
