# dev-libs/cyrus-sasl/files/saslauthd.rc6
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/dev-libs/cyrus-sasl/files/saslauthd.rc6,v 1.2 2003/02/14 22:40:59 vapier Exp $

depend() {
	need net
}

start() {
	ebegin "Starting saslauthd"
	start-stop-daemon --start --quiet --exec /usr/sbin/saslauthd \
		-- ${SASLAUTHD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping saslauthd"
	start-stop-daemon --stop --quiet --pidfile /var/lib/sasl/mux.pid
	eend $?
}
