# dev-libs/cyrus-sasl/files/saslauthd.rc6
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-libs/cyrus-sasl/files/saslauthd.rc6,v 1.6 2005/10/12 16:11:02 strerror Exp $

depend() {
	need net
	use ldap slapd
}

start() {
	ebegin "Starting saslauthd"
	start-stop-daemon --start --quiet --exec /usr/sbin/saslauthd \
		-- ${SASLAUTHD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping saslauthd"
	start-stop-daemon --stop --quiet --pidfile /var/lib/sasl/saslauthd.pid
	eend $?
}
