# dev-util/perforce-proxy/files/perforce-proxy-2003.1/init.d/perforce-proxy
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-util/perforce-proxy/files/perforce-proxy-2003.1/init.d/perforce-proxy,v 1.2 2004/03/04 18:44:32 vapier Exp $
#
# NB: Configuration is in /etc/conf.d/perforce-proxy

PIDFILE=/var/run/p4d.pid

depend() {
    need net
}


start() {
    ebegin "Starting Perforce SCM Proxy"
    # make sure they can't break our startup by passing -d (which would
    # cause incorrect PID to be written to pidfile), or the incorrect
    # depot root
    P4P="`echo ${P4P} | sed -e 's#-d##' -e 's#-r [^[:space:]]*##'`"
    start-stop-daemon --start --quiet --chuid perforce:perforce \
                      --exec /usr/sbin/p4p \
                      -- ${P4P} -d 
    eend $?

}

stop() {
    ebegin "Stopping Perforce SCM Proxy"
      if (  ps -C p4d | grep 'p4p' > /dev/null )
	  then
            killall p4d
	fi
    eend $?
}
