# sys-apps/pcsc-lite/files/pcscd-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-apps/pcsc-lite/files/pcscd-init,v 1.2 2004/03/06 04:09:07 vapier Exp $

start() {
	ebegin "Starting pcscd"
	start-stop-daemon --start --quiet --background --exec /usr/sbin/pcscd
	eend $?
}

stop() {
	ebegin "Stopping pcscd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/pcscd
	eend $?
}
