# media-gfx/maya/files/maya-docs
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting Maya Doc Server"
	start-stop-daemon --start --quiet --exec /opt/aw/maya/docs/startDocServer.sh
	eend $?
}

stop() {
	ebegin "Stopping Maya Doc Server"
	start-stop-daemon --start --quiet --exec /opt/aw/maya/docs/shutdownDocServer.sh
	eend $?
}
