# net-irc/ircd-hybrid/files/init.d_ircd-7.1.0
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-irc/ircd-hybrid/files/init.d_ircd-7.1.0,v 1.1 2005/07/20 19:11:02 swegener Exp $

opts="${opts} rehash"

depend() {
	need net
}

start() {
	ebegin "Starting IRCd-Hybrid7 daemon"
	start-stop-daemon --start --quiet --chuid hybrid --exec /usr/bin/ircd
	eend $?
}

stop() {
	ebegin "Stopping IRCd-Hybrid7 daemon"
	start-stop-daemon --stop --quiet --pidfile /var/run/ircd/ircd.pid
	eend $?
}

rehash() {
	ebegin "Rehashing IRCd-Hybrid7 IPv4 daemon"
	kill -s SIGHUP $(</var/run/ircd/ircd.pid)
	eend $?
}
