# net-nds/openldap/files/2.0/slurpd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-nds/openldap/files/2.0/slurpd,v 1.4 2004/07/15 00:18:02 agriffis Exp $

depend() {
	need net
}

start() {
	ebegin "Starting slurpd"
	start-stop-daemon --start --quiet \
	--exec /usr/lib/openldap/slurpd 
	eend $?
}

stop() {
	ebegin "Stopping slurpd"
	start-stop-daemon --stop --quiet \
	--exec /usr/lib/openldap/slurpd 
	eend $?
}
