# net-nds/portmap/files/portmap-5b-r6
#!/sbin/runscript

SERVICE=portmap
EXE=/sbin/portmap
opts="start stop restart"

depend() {
    need net
}

start() {
    ebegin "Starting ${SERVICE}"
    start-stop-daemon --start --quiet --exec $EXE 1>&2
    eend $? "Error starting ${SERVICE}."
}

stop() {
    ebegin "Stopping ${SERVICE}"
    start-stop-daemon --stop --quiet --exec $EXE 1>&2
    eend $? "Error stopping ${SERVICE}."
}

restart() {
    # Dump the portmapper's table before stopping
    ebegin "Saving portmap table"
    local tmpfile=`mktemp /tmp/pmap_table.XXXXXX`
    [ -n "$tmpfile" ] && pmap_dump >$tmpfile
    eend $? "Error saving portmap table."
    # Stop and restart portmapper
    stop
    start
    # Reload the portmapper's table
    if [ -n "$tmpfile" ]; then
	ebegin "Reloading portmap table"
	pmap_set <$tmpfile
	eend $? "Error reloading portmap table."
	rm -f $tmpfile
    fi
}
