# sys-apps/evms/files/evms2-init
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/sys-apps/evms/files/evms2-init,v 1.1 2003/04/07 08:05:54 tantive Exp $

start() {
	ebegin "Starting evms"
	evms_activate
	eend 0
}

stop() {
	ebegin "Stopping evms"
	evms_deactivate
	eend 0
}
