# net-mail/vm-pop3d/files/vm-pop3d.rc3
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-mail/vm-pop3d/files/vm-pop3d.rc3,v 1.6 2005/02/16 19:24:52 ferdy Exp $

depend() {
	need net
	use logger
}

start() {
	ebegin "Starting vm-pop3d"
	start-stop-daemon --start --quiet --exec /usr/sbin/vm-pop3d -- ${OPTS}
	eend ${?}
}

stop() {
	ebegin "Stopping vm-pop3d"
	killall vm-pop3d
	eend ${?}
}
