# net-mail/vm-pop3d/files/vm-pop3d.rc3
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-mail/vm-pop3d/files/vm-pop3d.rc3,v 1.5 2004/07/14 23:52:56 agriffis Exp $

SERVICE=vm-pop3d

opts="start stop restart"

start() {
	ebegin "Starting $SERVICE"
	/usr/sbin/vm-pop3d $OPTS
}

stop() {
	ebegin "Stopping $SERVICE"
	killall -15 /usr/sbin/vm-pop3d
}

restart() {
	stop
	start
}
