# net-mail/vm-pop3d/files/vm-pop3d.rc3
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/net-mail/vm-pop3d/files/vm-pop3d.rc3,v 1.2 2002/09/25 17:00:56 raker Exp $

SERVICE=vm-pop3d

opts="start stop restart"


start() {
		ebegin "Starting $SERVICE..."
     		/usr/sbin/vm-pop3d $OPTS
}

stop() {
		ebegin "Stopping $SERVICE..."
		killall -15 /usr/sbin/vm-pop3d
}

restart() {
		stop
		start
}

