# net-mail/postfix/files/postfix
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/net-mail/postfix/files/postfix,v 1.4 2003/04/19 02:55:27 drobbins Exp $

depend() {
	need net
	use logger dns
	provide mta
}

PIDFILE=/var/spool/postfix/pid/master.pid

start() {
	ebegin "Starting postfix"
	/usr/sbin/postfix start &>/dev/null
	eend $?
}

stop() {
	ebegin "Stopping postfix"
	/usr/sbin/postfix stop &>/dev/null
	eend $?
}
