# net-mail/postfix/files/postfix-rc6
#!/sbin/runscript

PIDFILE=/var/spool/postfix/pid/master.pid

depend() {
    need sysklogd net
}

start() {
    ebegin "Starting postfix"
    /usr/sbin/postfix start >/dev/null 2>&1
    eend $?
}

stop() {
    ebegin "Stopping postfix"
    /usr/sbin/postfix stop >/dev/null 2>&1
    eend $?
}
