# net-mail/postfix/files/postfix-rc6
#!/sbin/runscript

SERVICE=postfix
opts="start stop restart"

PIDFILE=/var/spool/postfix/pid/master.pid
EXE=/usr/sbin/postfix

depend() {
        need sysklogd net
}

start() {
        ebegin "Starting $SERVICE..."
        $EXE start >/dev/null 2>&1
        eend $? "Started $SERVICE." "Error Starting $SERVICE."
}

stop() {
        ebegin "Stopping $SERVICE..."
        $EXE stop >/dev/null 2>&1
        eend $? "Stopped $SERVICE." "Error Stopping $SERVICE."
}

restart() {
        stop
        start
}
