# net-mail/postfix/files/postfix.rc6
#!/sbin/runscript

depend() {
    need net
    use sysklogd metalog syslog-ng
}


PIDFILE=/var/spool/postfix/pid/master.pid

start() {
    ebegin "Starting postfix"
    /usr/sbin/postfix start >/dev/null 2>&1
    eend $?
}

stop() {
    ebegin "Stopping postfix"
    /usr/sbin/postfix stop >/dev/null 2>&1
    eend $?
}
