# net-www/memcached/files/init
#!/sbin/runscript
# $Header: /var/cvsroot/gentoo-x86/net-www/memcached/files/init,v 1.3 2003/10/09 18:29:25 lisa Exp $

depend() {
	need net
}

start() {
	ebegin "Starting memcached"
	/sbin/start-stop-daemon --start \
	        --exec ${MEMCACHED_BINARY} --pidfile /var/run/memcached.pid \
		-- -d -p ${PORT} -l ${LISTENON} -m ${MEMUSAGE} \
		-c ${MAXCONN} -u ${MEMCACHED_RUNAS}
	eend $?
}

stop() {
	ebegin "Stopping memcached"
	start-stop-daemon --stop --quiet --exec ${MEMCACHED_BINARY}
	rm -f /var/run/memcached.pid
	eend $?
}

