# net-analyzer/nessus-core/files/nessusd-r6
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-analyzer/nessus-core/files/nessusd-r6,v 1.4 2004/07/14 23:05:07 agriffis Exp $

depend() {
	need net
}

start() {
	ebegin "Starting nessusd"
	start-stop-daemon --start --quiet --exec /usr/sbin/nessusd -- -D
	eend $? 
}

stop() {
	ebegin "Stop nessusd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/nessusd
	eend $?
}
