# sys-apps/daemontools/files/svscan
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/sys-apps/daemontools/files/svscan,v 1.4 2002/07/27 00:41:15 raker Exp $

depend() {
    need net
}

start() {
    ebegin "Starting service scan..."
    PATH="/usr/bin:$PATH"
    start-stop-daemon --start --quiet --exec /usr/bin/svscan --pidfile /var/run/svscan.pid -- /service &
    echo $! > /var/run/svscan.pid
    eend $?
}

stop() {
    ebegin "Stopping service scan..."
    start-stop-daemon --stop --quiet --pidfile /var/run/svscan.pid
    eend $?
    ebegin "Stopping services..."
    svc -dx /service/*
    eend $?
    ebegin "Stopping service logging..."
    svc -dx /service/*/log
    eend $?
}

