# net-wireless/bluez-utils/files/bluez-utils-3.10.1-udev.rules
# Start/Stop bluetooth service on device insertion. Gentoo specific.
SUBSYSTEM=="bluetooth", KERNEL=="hci[0-9]*", RUN+="bluetooth.sh"
