# net-fs/coda/files/auth2
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-fs/coda/files/auth2,v 1.1 2004/03/23 20:18:49 dragonheart Exp $

depend() {
	need net
}

checkconfig() {
	if [ ! -e /var/lib/vice/hostname ] ; then
		eerror "Please set up vice before starting the service..."
		return 1
	fi
}

we_are_scm() {
	if [ "`cat /var/lib/vice/hostname`" = "`cat /var/lib/vice/db/scm`" ] ; then
		return 0
	else
		return 1
	fi
}

start() {
	checkconfig || return 1

	ebegin "Starting auth2"

	# Check to see if we are the SCM. If not, we must add -chk to auth2

	if (! we_are_scm) ; then
		start-stop-daemon --start --quiet --exec /usr/sbin/auth2 \
			--background -- -chk
	else
		start-stop-daemon --start --quiet --exec /usr/sbin/auth2 \
			--background
	fi

	eend $?
}

stop() {
	ebegin "Stopping auth2"
	start-stop-daemon --stop --quiet --exec /usr/sbin/auth2
	eend $?
}
