# net-fs/coda/files/codasrv
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-fs/coda/files/codasrv,v 1.3 2004/07/16 23:22:03 dragonheart Exp $

depend() {
	need net coda-update auth2
}

checkconfig() {
        if [ ! -e /etc/coda/server.conf ] ; then
                eerror "Please run vice-setup before starting the service..."
                return 1
        fi
        eval `grep "^vicedir=" /etc/coda/server.conf || echo vicedir=/var/lib/vice`
        if [ ! -e $vicedir/hostname ] ; then
                eerror "Please run vice-setup before starting the service..."
                return 1
        fi
}

we_are_scm() {
	if [ "`cat $vicedir/hostname`" = "`cat $vicedir/db/scm`" ] ; then
		return 0
	else
		return 1
	fi
}

start() {
	checkconfig || return 1

	ebegin "Starting codasrv"
	start-stop-daemon --start --exec /usr/sbin/codasrv
	eend $?
}

stop() {
	ebegin "Stopping codasrv"
	start-stop-daemon --stop --quiet --exec /usr/sbin/codasrv
	/usr/sbin/volutil shutdown
	sleep 30
	eend $?
}
