# net-fs/coda/files/venus
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-fs/coda/files/venus,v 1.1 2004/03/23 20:18:49 dragonheart Exp $

depend() {
	need net codasrv
}

checkconfig() {
	if [ ! -e /var/lib/vice/hostname ] ; then
		eerror "Please set up vice before starting the service..."
		return 1
	fi
}

start() {
	checkconfig || return 1

	ebegin "Starting venus"
	start-stop-daemon --start --quiet --exec /usr/sbin/venus --background
	eend $?
}

stop() {
	ebegin "Stopping venus"
	start-stop-daemon --stop --quiet --exec /usr/sbin/venus
	umount -l /mnt/coda &>/dev/null
	eend $?
}
