# sys-devel/distcc/files/distccd.4
#!/sbin/runscript

DISTCCD_EXEC=/usr/bin/distccd
DISTCCD_PIDFILE=/var/run/distccd/distccd.pid

depend() {
	need net
}

start() {
	rm -f ${DISTCCD_PIDFILE}
	ebegin "Starting distccd"
	PATH="$(gcc-config --get-bin-path):${PATH}" \
	    /sbin/start-stop-daemon --start --quiet \
	        --startas ${DISTCCD_EXEC} \
		-c distcc --pidfile ${DISTCCD_PIDFILE} \
		-- --pid-file ${DISTCCD_PIDFILE}
	eend $? "Failed to Start distccd"
}

stop() {
	ebegin "Stopping distccd"
	start-stop-daemon --stop --quiet --pidfile "${DISTCCD_PIDFILE}"
	rm -f "${DISTCCD_PIDFILE}"
	eend $? "Failed to Stop distccd"
}

restart() {
	ebegin "Restarting distccd"
	svc_stop
	svc_start
	eend $? "Failed to Restart distccd"
}
