# sys-devel/distcc/files/distccd
#!/sbin/runscript

depend() {
	need net
}

start() {
	rm -rf /var/tmp/distccd.pid
	ebegin "Starting distccd"
	start-stop-daemon --start --quiet --startas /usr/bin/distccd \
		-c nobody --pidfile /var/tmp/distccd.pid \
		-- --pid-file /var/tmp/distccd.pid
	eend $? "Failed to Start distccd"
}

stop() {
	ebegin "Stopping distccd"
	start-stop-daemon --stop --quiet --pidfile /var/tmp/distccd.pid
	rm -rf /var/tmp/distccd.pid
	eend $? "Failed to Stop distccd"
}

restart() {
	ebegin "Restarting distccd"
	svc_stop
	svc_start
	eend $? "Failed to Restart distccd"
}
