# sys-devel/distcc/files/distccd.3
#!/sbin/runscript

depend() {
	need net
}

start() {
	rm -rf /var/run/distccd.pid
	ebegin "Starting distccd"
	PATH="$(gcc-config --get-bin-path):${PATH}" \
	    /sbin/start-stop-daemon --start --quiet \
	        --startas /usr/bin/distccd \
		-c nobody --pidfile /var/run/distccd.pid \
		-- --pid-file /var/run/distccd.pid
	eend $? "Failed to Start distccd"
}

stop() {
	ebegin "Stopping distccd"
	start-stop-daemon --stop --quiet --pidfile /var/run/distccd.pid
	rm -rf /var/run/distccd.pid
	eend $? "Failed to Stop distccd"
}

restart() {
	ebegin "Restarting distccd"
	svc_stop
	svc_start
	eend $? "Failed to Restart distccd"
}
