# net-irc/dancer-ircd/files/dancer-ircd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-irc/dancer-ircd/files/dancer-ircd,v 1.3 2004/03/06 03:11:35 vapier Exp $

depend() {
	need net
}

start() {
	ebegin "Starting dancer-ircd"
	start-stop-daemon --start --quiet --exec /usr/bin/dancer-ircd
	eend $?
}

stop() {
	ebegin "Stopping dancer-ircd"
	start-stop-daemon --stop --quiet --pidfile /var/run/dancer-ircd.pid
	eend $?
}
