# net-irc/dancer-ircd/files/dancer-ircd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/net-irc/dancer-ircd/files/dancer-ircd,v 1.2 2004/01/03 04:05:52 zul Exp $

depend() {
	need net
}

start() {
	ebegin "Starting dancer-ircd"
	start-stop-daemon --start --quiet --exec /usr/bin/dancer-ircd
	result=$?
	eend $result
}

stop() {
	ebegin "Stopping dancer-ircd"
	start-stop-daemon --stop --quiet --pidfile /var/run/dancer-ircd.pid
	result=$?
	eend $result
}
