# net-irc/dancer-ircd/files/dancer-ircd-1.0.35
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-irc/dancer-ircd/files/dancer-ircd-1.0.35,v 1.2 2005/07/31 20:41:03 swegener Exp $

depend() {
	use net
	provide ircd
}

start() {
	ebegin "Starting dancer-ircd"
	start-stop-daemon --start --quiet --chuid ${DANCER_USER} --exec /usr/bin/dancer-ircd
	eend $?
}

stop() {
	ebegin "Stopping dancer-ircd"
	start-stop-daemon --stop --quiet --pidfile /var/run/dancer-ircd/dancer-ircd.pid
	eend $?
}
