# www-misc/xxv/files/xxv
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/www-misc/xxv/files/xxv,v 1.1 2009/04/11 23:02:00 hd_brummy Exp $

RUNAS_USER="vdr"

# some fixed Path
CONFIGFILE="/etc/xxv/xxvd.cfg"
PIDFILE="/var/run/xxv/xxvd.pid"
XXV_BIN="/usr/bin/xxvd"

LANGUAGE="de_DE"

# Set Verbose Level 0 -> 5 
VERBOSE="3"

#set -xv

depend() {
	need vdr
	after net
}

kill_vdradmin() {

		# Check at first, is VDR-Admin running
		# Stopping, while running on same ports
		# You can comment the next 4 Lines if you running VDR-Admind on different Ports
		if [ -f /var/run/vdradmind.pid ] ; then
		ebegin "VDR-Admin will Stop at first now"
		/etc/init.d/vdradmind stop
		fi
}

xxv_kill_pid() {

		# After unclear stop, xxvd.pid will not removed, fixed with next line
		if [ -e ${PIDFILE} -a ! -L /var/lib/init.d/started/xxv ] ; then
			rm ${PIDFILE}
			killall xxvd 2>&1 > /dev/null
		fi
}

xxv_kill_initfile() {

		# After unclear stop, init file in /var/lib/init.d/started/ still not removed
		if [ -L /var/lib/init.d/started/xxv -a ! -e ${PIDFILE} ] ; then
			rm /var/lib/init.d/started/xxv
			/etc/init.d/xxv zap
			killall xxvd 2>&1 > /dev/null
		fi
}
	
start() {

	kill_vdradmin

	# Finaly start XXV :)
	export LANG="${LANGUAGE}"
	ebegin "Start xxv"
	start-stop-daemon  --nicelevel 15 --pidfile ${PIDFILE} --start -c ${RUNAS_USER} --exec ${XXV_BIN} -- \
	-configfile=${CONFIGFILE} -verbose=${VERBOSE} -pidfile=${PIDFILE}
	eend $?
}


stop() {

	ebegin "Stopping xxv"
	start-stop-daemon --stop --quiet --pidfile ${PIDFILE}

	xxv_kill_pid
	xxv_kill_initfile

	eend $?
}

restart() {

	xxv_kill_pid
	xxv_kill_initfile


	svc_stop
	svc_start
}
