# sys-apps/man/files/makewhatis.cron
#!/bin/sh
#
# Uncomment the following if you want 
# your whatis database rebuilt regularly
#
# this is part of the man package
#
#

/usr/sbin/makewhatis
