# net-firewall/shorewall/files/shorewall
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-firewall/shorewall/files/shorewall,v 1.4 2004/03/04 19:20:40 vapier Exp $

opts="start stop restart"

depend() {
	need net
	provide firewall
}

start() {
	ebegin "Starting firewall"
	/sbin/shorewall start 1>/dev/null
	eend $? 
}

stop() {
	ebegin "Stopping firewall"
	/sbin/shorewall stop 1>/dev/null
	eend $?
}

restart() {
	ebegin "Restarting firewall"
	/sbin/shorewall restart
	eend $?
}
