# app-admin/wasabi/files/wasabi.init
#!/sbin/runscript

depend() {
	need clock hostname logger
}

start() {
	ebegin "Starting wasabi"
	start-stop-daemon --start --quiet --exec /usr/sbin/wasabi -- -c /etc/wasabi/wasabi.conf
	eend $?
}

stop() {
	ebegin "Stopping wasabi"
	start-stop-daemon --stop --quiet --name wasabi
	eend $?
}
