# app-emulation/win4lin/files/win4lin.initd.new
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-emulation/win4lin/files/win4lin.initd.new,v 1.2 2004/02/13 19:13:59 vapier Exp $

depend() {
	need net
}

start() {
	ebegin "Starting Win4Lin"
	/etc/rc.d/init.d/Win4Lin start > /dev/null
	eend $?
}

stop() {
	ebegin "Stopping Win4Lin"
	/etc/rc.d/init.d/Win4Lin stop > /dev/null
	eend $?
}

