# sys-apps/apcupsd/files/apcupsd
#!/sbin/runscript
 
SERVICE=apcupsd
EXE=/usr/sbin/apcupsd
 
start() {
	rm -f /etc/apcupsd/powerfail
	ebegin "Starting ${EXE}"
	start-stop-daemon --start --quiet --exec ${EXE} -- 1>&2
	eend $? "Error starting ${EXE}."
}
 
stop() {
	ebegin "Stopping ${EXE}"
	start-stop-daemon --stop --quiet --exec ${EXE} -- 1>&2
	eend $? "Error stopping ${EXE}."
} 
