# sys-apps/apcupsd/files/apcupsd
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/sys-apps/apcupsd/files/apcupsd,v 1.7 2003/06/29 22:38:10 mholzer Exp $

depend() {
	need net
	use logger dns
}

 
start() {
	rm -f /etc/apcupsd/powerfail
	ebegin "Starting APC UPS daemon"
	start-stop-daemon --start --quiet --exec /usr/sbin/apcupsd -- 1>&2
	eend $?
}
 
stop() {
	ebegin "Stopping APC UPS daemon"
	start-stop-daemon --stop --quiet --exec /usr/sbin/apcupsd -- 1>&2
	eend $?
} 

restart() {
	svc_stop
	# need to sleep or new process may choke grabbing
	#  the dying processes network socket
	sleep 1
	svc_start
}
