# net-p2p/freenet/files/freenet
#!/sbin/runscript
# Freenet init.d-script by Per Wigren <wigren@home.se>


depend () {
	need net
}

get_seednodes () {
	if [[ ! -e /var/freenet/seednodes.ref ]]; then
		echo No seednodes.ref file was found. Fetching a new one...
		cd /var/freenet
		wget -q http://hawk.freenetproject.org/~freenet4/seednodes.ref
	fi
}

start () {
	ebegin "Starting Freenet"
	get_seednodes
	export CLASSPATH=/usr/lib/freenet/freenet.jar:/usr/lib/freenet/freenet-ext.jar:$CLASSPATH
	start-stop-daemon --start --quiet --pidfile /var/run/freenet.pid -m \
		--background \
		--exec $(java-config --java) -- freenet.node.Main -p /etc/freenet.conf
	eend $?
}

stop () {
	ebegin "Stopping Freenet"
	start-stop-daemon --stop --quiet --pidfile /var/run/freenet.pid
	rm -f /var/run/freenet.pid
	eend $?
}

