# net-p2p/freenet/files/freenet
#!/sbin/runscript
# Freenet init.d-script by Per Wigren <wigren@home.se>

depend() {
	need net
}

SEEDNODES_REF="http://freenetproject.org/snapshots/seednodes.ref"
FREENET_JAR="http://freenetproject.org/snapshots/freenet-latest.jar"

start() {
	einfo "Fetching latest seednodes.ref..."
	mv -f /var/freenet/seednodes.ref /var/freenet/seednodes.ref.old &>/dev/null
	wget -O /var/freenet/seednodes.ref -q $SEEDNODES_REF || \
		mv -f /var/freenet/seednodes.ref.old /var/freenet/seednodes.ref

	einfo "Fetching latest freenet.jar..."
	mv -f /usr/lib/freenet/freenet.jar /usr/lib/freenet/freenet.jar.old &>/dev/null
	wget -O /usr/lib/freenet/freenet.jar -q $FREENET_JAR || \
		mv -f /usr/lib/freenet/freenet.jar.old /usr/lib/freenet/freenet.jar

	ebegin "Starting Freenet"
	export CLASSPATH=/usr/lib/freenet/freenet.jar:/usr/lib/freenet/freenet-ext.jar:$CLASSPATH
	start-stop-daemon --start --quiet --pidfile /var/run/freenet.pid -m \
		--background \
		--exec $(java-config --java) -- freenet.node.Main -p /etc/freenet.conf
	eend $?
}

stop() {
	ebegin "Stopping Freenet"
	start-stop-daemon --stop --quiet --pidfile /var/run/freenet.pid
	rm -f /var/run/freenet.pid
	eend $?
}
