# app-forensics/lynis/files/lynis.cron
#!/bin/sh

/usr/bin/lynis -c --auditor "automated" --cronjob
