# sys-apps/hdparm/files/hdparm-init-2
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Copyright 2003 Justin "Nagash" Jones
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-apps/hdparm/files/hdparm-init-2,v 1.1 2004/08/01 14:43:43 lanius Exp $

depend() {
	need localmount
}

checkconfig() {
	drives=""
	cdroms=""

	for i in `dir /dev/discs 2>/dev/null`
	do
		if [ ! -n "`ls -la /dev/discs/${i} | grep scsi`" ]
		then
			all_discs="${all_discs} `basename ${i}`"
		elif [ -n "${all_args}" ]
		then
			ewarn "Skipping $i hdparm does not support SCSI devices."
		fi
	done
	for i in `dir /dev/cdroms 2>/dev/null`
	do	
		if [ ! -n "`ls -la /dev/cdroms/${i} | grep scsi`" ]
		then
			all_cdroms="${all_cdroms} `basename ${i}`"
		elif [ -n "${all_args}" ]
		then
			ewarn "Skipping $i hdparm does not support SCSI devices."
		fi
	done

	for i in ${all_discs}
	do
		tmp_args=${i}_args
		eval hd_args=\$${tmp_args}

		if [ -n "${hd_args}" ]
		then
			drives="$drives $i"
		fi
	done

	for i in ${all_cdroms}
	do
		tmp_args=${i}_args
		eval hd_args=\$${tmp_args}

		if [ -n "${hd_args}" ]
		then
			cdroms="$cdroms $i"
		fi
	done
}

start() {
	ebegin "Starting hdparm"
	checkconfig
	
	if [ -n "${all_args}" ]
	then	
		for i in ${all_discs}
		do
			ebegin "Running hdparm on ${i}"
			hdparm ${all_args} /dev/discs/${i}/disc &> /dev/null
			eend $? "Failed to start hdparm on ${i}."
		done

		for i in ${all_cdroms}
		do
			ebegin "Running hdparm on ${i}"
			hdparm ${all_args} /dev/cdroms/${i} &> /dev/null
			eend $? "Failed to start hdparm on ${i}."
		done
	else
		for i in ${drives}
		do
			tmp_args=${i}_args
			eval hd_args=\$$tmp_args
			ebegin "Running hdparm on ${i}"
			hdparm $hd_args /dev/discs/${i}/disc &> /dev/null
			eend $? "Failed to start hdparm on ${i}."
		done

		for i in ${cdroms}
		do
			tmp_args=${i}_args
			eval hd_args=\$$tmp_args
			ebegin "Running hdparm on ${i}"
			hdparm $hd_args /dev/cdroms/${i} &> /dev/null
			eend $? "Failed to start hdparm on ${i}."
		done
	fi

	eend $? "Failed to start hdparm."
}

stop() {
	ebegin "Stopping hdparm"
	eend $? "Failed to stop hdparm."
}
