# sys-apps/hdparm/files/hdparm-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-apps/hdparm/files/hdparm-init,v 1.4 2004/07/15 00:43:00 agriffis Exp $

depend() {
	need localmount
}

start() {
	ebegin "Enabling DMA on IDE drives"
	/sbin/hdparm -d1 /dev/ide/hd/*u? &> /dev/null
	eend $?
}
