# app-emulation/virtualbox-ose/files/vboxwebsrv-initd
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-emulation/virtualbox-ose/files/vboxwebsrv-initd,v 1.3 2009/11/21 12:01:15 patrick Exp $

depend() {
	need net
}

start() {
	#checkconfig || return $?

	ebegin "Starting VirtualBox WebService"
	start-stop-daemon --start --quiet --pidfile /var/run/vboxwebsrv.pid \
		--background --make-pidfile --exec /opt/bin/vboxwebsrv \
		-- -H ${VBOXWEBSRV_HOST} -p ${VBOXWEBSRV_PORT} -t ${VBOXWEBSRV_TIMEOUT} -i ${VBOXWEBSRV_INTERVAL}
	eend $?
}

stop() {
        ebegin "Stopping VirtualBox WebService"
        start-stop-daemon --stop --quiet --pidfile /var/run/vboxwebsrv.pid
        eend $?
}
