# net-news/inn/files/innd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-news/inn/files/innd,v 1.7 2004/09/07 22:29:04 swegener Exp $

depend() {
	need net
}

start() {
	ebegin "Starting innd"
	su - news -c /usr/lib/news/bin/rc.news
	eend $?
}

stop() {
	ebegin "Stopping innd"
	su - news -c '/usr/lib/news/bin/rc.news stop'
	sleep 2
	eend $?
}
