# net-fs/samba/files/samba-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-fs/samba/files/samba-init,v 1.4 2004/03/04 19:25:50 vapier Exp $

opts="${opts} reload"

depend() {
	need net
	use cupsd
}

start() {
	ebegin "Starting samba"
	start-stop-daemon --start --quiet --exec /usr/sbin/smbd
	result=$?
	start-stop-daemon --start --quiet --exec /usr/sbin/nmbd
	result=$(( $result + $? ))
	eend $result
}

stop() {
	ebegin "Stopping samba"
	start-stop-daemon --stop --quiet --pidfile /var/run/samba/smbd.pid
	result=$?
	start-stop-daemon --stop --quiet --pidfile /var/run/samba/nmbd.pid
	result=$(( $result + $? ))
	eend $result
}

reload() {
	ebegin "Reloading samba"
	killall -HUP smbd &>/dev/null
	result=$?
	killall -HUP nmbd &>/dev/null
	result=$(( $result + $? ))
	eend $result
}
