# net-fs/samba/files/winbind-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-fs/samba/files/winbind-init,v 1.5 2004/07/14 23:37:16 agriffis Exp $

depend() {
	need net
}

start() {
	ebegin "Starting winbind"
	start-stop-daemon --start --quiet --exec /usr/sbin/winbindd
	eend $?
}

stop() {
	ebegin "Stopping winbind"
	start-stop-daemon --stop --quiet --oknodo --exec /usr/sbin/winbindd
	eend $?
}
