# media-sound/alsa-utils/files/alsa-0.5.10
#!/sbin/runscript

depend() {
	need modules
}

start() {
	if [ -e /etc/asound.conf ]
	then
		ebegin "Restoring ALSA volume settings"
		/usr/sbin/alsactl restore
	else
		ebegin "/etc/asound.conf not found, ALSA muted"
	fi
    eend $? 
}

stop() {
	ebegin "Saving ALSA volume settings"
	/usr/sbin/alsactl store
	eend 
	#will fail if no modules were loaded, so we don't process the return val
}
