# net-www/vdradmin-am/files/vdradmin-4
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/net-www/vdradmin-am/files/vdradmin-4,v 1.2 2009/01/13 18:00:03 hd_brummy Exp $

depend() {
	use vdr
}

set_charset() {
	local capfile=/usr/share/vdr/capabilities.sh
	[ -e "${capfile}" ] && . ${capfile}
	if [ "${CAP_UTF8}" = "1" -o "${VDR_CAN_HANDLE_UTF8}" = "yes" ]; then
		# do not clean out utf8
		:
	else
		l=$(locale | sed -e 's/\.[uU][tT][fF]\(-\)\?8//' -e 's/^/export /')
		eval $l
	fi
}

start() {
	ebegin "Starting vdradmin-am"

	: ${VDRADMIN_USER:=vdradmin}

	# Guarantee that vdradmin can access its tmp-dir
	chown -R "${VDRADMIN_USER}:vdradmin" /var/tmp/vdradmin

	set_charset
	if [ -x /etc/init.d/root ]; then
		# running on baselayout-2/openrc
		start-stop-daemon --start --background --chuid "${VDRADMIN_USER}" --exec /usr/bin/perl /usr/bin/vdradmind
	else
		start-stop-daemon --start --background --chuid "${VDRADMIN_USER}" --exec /usr/bin/vdradmind
	fi
	eend $?
}

stop() {
	ebegin "Stopping vdradmin-am"
	start-stop-daemon --stop --quiet --pidfile /var/tmp/vdradmin/vdradmind.pid
	eend $?
}
