# sys-apps/anacron/files/anacron.rc6
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-apps/anacron/files/anacron.rc6,v 1.4 2004/07/15 00:36:54 agriffis Exp $

depend() {
	need clock hostname logger
}

start() {
	ebegin "Running anacron"
	/usr/sbin/anacron -s >>/var/log/cron.log 2>&1
	eend $?
}
