# sys-apps/anacron/files/anacron.rc6
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-apps/anacron/files/anacron.rc6,v 1.5 2004/10/08 16:30:26 ka0ttic Exp $

depend() {
	use logger
	need clock hostname 
}

start() {
	ebegin "Running anacron"
	/usr/sbin/anacron -s >>/var/log/cron.log 2>&1
	eend $?
}
