# sys-apps/anacron/files/anacron.rc6
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/sys-apps/anacron/files/anacron.rc6,v 1.2 2003/02/14 23:11:46 vapier Exp $

depend() {
	need clock hostname logger
}

start() {
	ebegin "Running anacron"
	/usr/sbin/anacron -s >>/var/log/cron.log 2>&1
	eend $?
}
