# net-fs/sfs/files/sfssd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-fs/sfs/files/sfssd,v 1.3 2004/07/14 23:39:36 agriffis Exp $

depend() {
	need net portmap nfs
}

start() {
	ebegin "Starting SFS server daemon"
	start-stop-daemon --start --quiet --exec /sbin/sfssd 1>&2
	eend $?
}

stop() {
	# Stopping/restarting will likely cause any connected SFS clients
	# to hang messily until they are restarted.	 Just keep this in
	# mind, as I do not believe that the current version has a
	# work-around.

	ebegin "Stopping SFS server daemon"
	start-stop-daemon --stop --quiet --exec /sbin/sfssd 1>&2
	eend $?
}
