# app-antivirus/vlnx/files/uvscan.cron
#!/bin/sh

/opt/vlnx/virusscan_updater.sh
/opt/bin/uvscan --config /etc/uvscan.conf /home
