# app-antivirus/vlnx/files/uvscan.cron
#!/bin/bash

/opt/vlnx/check-updates.sh
uvscan --config /etc/uvscan.conf /home
