# net-nntp/inn/files/innd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-nntp/inn/files/innd,v 1.3 2012/02/08 06:41:07 slyfox Exp $

depend() {
	need net
}

start() {
	ebegin "Starting innd"
	start-stop-daemon --start --user news --exec /usr/lib/news/bin/rc.news
	eend $?
}

stop() {
	ebegin "Stopping innd"
	start-stop-daemon --start --user news --exec /usr/lib/news/bin/rc.news -- stop
	sleep 2
	eend $?
}
