# app-admin/watchdog/files/5.2/watchdog
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/watchdog/files/5.2/watchdog,v 1.5 2004/07/14 21:21:38 agriffis Exp $

depend() {
	need localmount
	use logger
}

checkconfig() {
	if [ ! -e /dev/watchdog ] ; then
		eerror "/dev/watchdog does not exist!"
		return 1
	fi
	return 0
}

start() {
	checkconfig || return 1
	ebegin "Starting watchdog"
	start-stop-daemon --start --quiet --exec /usr/sbin/watchdog
	eend $?
}

stop() {
	ebegin "Stopping watchdog"
	start-stop-daemon --stop --quiet --exec /usr/sbin/watchdog
	eend $?
}
