# net-misc/sobby/files/sobby-init-0.2.0
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2         
# $Header: /var/cvsroot/gentoo-x86/net-misc/sobby/files/sobby-init-0.2.0,v 1.1 2005/08/14 22:01:58 humpback Exp $

depend() {
	need mDNSResponder
}

start() {
	ebegin "Starting obby dedicated server"
	start-stop-daemon -b --quiet --start --exec /usr/bin/sobby -- ${SOBBY_OPTS}
	eend $? "Failed to start sobby"
}

stop() {
	ebegin "Stopping obby dedicated server"
	start-stop-daemon --stop --quiet --exec /usr/bin/sobby
	eend $? "Failed to stop sobby"
}
