# net-misc/ssh/files/sshd2
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/ssh/files/sshd2,v 1.3 2004/07/15 00:12:48 agriffis Exp $

depend() {
	use logger dns
	need net
}

checkconfig() {
	if [ ! -e /etc/ssh2/sshd2_config ] ; then
		eerror "You need an /etc/ssh/sshd2_config file to run sshd"
		eerror "There is a sample file in  /etc/ssh2/"
		return 1
	fi
}


start() {
	checkconfig || return 1
	ebegin "Starting sshd2"
	start-stop-daemon --start --quiet --pidfile /var/run/sshd2.pid \
		--startas /usr/sbin/sshd2
	eend $?
}

stop() {
	ebegin "Stopping sshd2"
	start-stop-daemon --stop --quiet --pidfile /var/run/sshd2_22.pid
	eend $?
}
