# app-text/man2html/files/man2html.cron
#!/bin/sh
nice -n 19 /usr/bin/glimpseindex -f -z -H /var/man2html `man -w | tr : ' '`
chmod 644 /var/man2html/.glimpse*
