# sys-libs/libchipcard/files/chipcardd3
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-libs/libchipcard/files/chipcardd3,v 1.2 2007/01/27 11:16:43 hanno Exp $

start() {
    ebegin "Starting chipcard3 daemon"
    start-stop-daemon --start --quiet --exec /usr/sbin/chipcardd3
    eend $?
}

stop() {
    ebegin "Stopping chipcard3 daemon"
    start-stop-daemon --stop --quiet --exec /usr/sbin/chipcardd3
    eend $?
}
