# sys-libs/libchipcard/files/chipcardd4
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-libs/libchipcard/files/chipcardd4,v 1.3 2010/01/07 10:40:19 ssuominen Exp $

depend() {
    need hald
}

start() {
    ebegin "Starting chipcard4 daemon"
    start-stop-daemon --start --quiet --exec /usr/sbin/chipcardd4
    eend $?
}

stop() {
    ebegin "Stopping chipcard4 daemon"
    start-stop-daemon --stop --quiet --exec /usr/sbin/chipcardd4
    eend $?
}
