# net-dialup/pptpclient/files/pptp-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-dialup/pptpclient/files/pptp-init,v 1.2 2004/09/28 01:22:04 swegener Exp $

depend() {
        need net
}

start() {

        ebegin "Starting Point-to-Point Tunnelling Server"
        start-stop-daemon --start --quiet --background --pidfile /var/run/pptpd.pid --make-pidfile --exec /usr/sbin/pptpd
        eend $?
}

stop() {
        ebegin "Stopping Point-to-Point Tunnelling Server"
        start-stop-daemon --stop --quiet --pidfile /var/run/pptpd.pid --name pptpd
        eend $?
}
