# www-proxy/squid/files/squid.cron
#!/bin/sh
exec /usr/sbin/squid -k rotate
