# app-misc/uptimed/files/uptimed
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-misc/uptimed/files/uptimed,v 1.4 2004/07/14 21:34:49 agriffis Exp $

# Copyright 2002 Sean Russell

start() {
	ebegin "Starting uptimed"
	uptimed -b
	start-stop-daemon --quiet --start --pidfile /var/run/uptimed --exec /usr/sbin/uptimed
	eend $?
}

stop() {
	ebegin "Stopping uptimed"
	start-stop-daemon --quiet --stop --pidfile /var/run/uptimed
	eend $?
}
