# net-analyzer/nagios-nsca/files/nsca.init
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-analyzer/nagios-nsca/files/nsca.init,v 1.2 2012/08/21 17:30:54 flameeyes Exp $

get_pidfile() {
    sed -n -e 's:^[ \t]*pid_file=\([^#]\+\).*:\1:p' \
	/etc/nagios/nsca.cfg
}

extra_started_commands="reload"

command=/usr/libexec/${SVCNAME}
command_args="-c /etc/nagios/nsca.cfg --daemon"
pidfile=$(get_pidfile)

depend() {
    need net
    use icinga nagios
}

reload() {
    ebegin "Reloading nsca"
    kill -HUP `cat $(get_pidfile)`
    eend $? "Failed to reload nsca"
}
