# net-dialup/l2tpd/files/l2tpd-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: 

depend() {
	need net
}

start() {
	ebegin "Starting l2tpd"
	start-stop-daemon --start --quiet --exec /usr/sbin/l2tpd
	eend $?
}

stop() {
	ebegin "Stopping l2tpd"
	start-stop-daemon --stop --quiet --pidfile /var/run/l2tpd.pid
	eend $?
}
