# sys-fs/evms/files/evms2-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/sys-fs/evms/files/evms2-init,v 1.2 2004/01/24 08:37:30 seemant Exp $

start() {
	ebegin "Starting evms"
	evms_activate
	eend $?
}

stop() {
	ebegin "Stopping evms"
	eend 0
}
