# www-servers/tomcat/files/5.0.28/tomcat.init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/www-servers/tomcat/files/5.0.28/tomcat.init,v 1.5 2005/09/22 00:06:03 betelgeuse Exp $

depend() {
	use dns logger mysql postgresql net
}

start()	{
	ebegin "Starting Tomcat"
	cd ${CATALINA_BASE}/webapps
	start-stop-daemon --start --quiet --chuid tomcat:tomcat \
			  --exec ${CATALINA_HOME}/bin/catalina.sh \
			  --background -- ${TOMCAT_START}
	#sleep 5
	eend $?
}

stop()	{
	ebegin "Stopping Tomcat"
	start-stop-daemon --start --quiet --chuid tomcat:tomcat \
			  --exec ${CATALINA_HOME}/bin/catalina.sh \
		 	  --background -- ${TOMCAT_STOP}
	#sleep 5
	eend $?
}
