# mail-mta/qmail-mysql/files/1.03/qmail-control
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/mail-mta/qmail-mysql/files/1.03/qmail-control,v 1.2 2004/07/14 22:34:55 agriffis Exp $

opts="start stop restart reload"

start() {
    ebegin "Starting qmail mta ..."
    svc -u /var/qmail/supervise/*
    eend $?
    ebegin "Starting qmail mta logging ..."
    svc -u /var/qmail/supervise/*/log
    eend $?
}

stop() {
    ebegin "Stopping qmail mta ..."
    svc -d /var/qmail/supervise/*
    eend $?
    ebegin "Stopping qmail mta logging ..."
    svc -d /var/qmail/supervise/*/log
    eend $?
}

reload() {
    ebegin "Reloading 'locals' and 'virtualdomains' control files."
    svc -h /var/qmail/supervise/qmail-send
    eend $?
}
