# www-apps/nut/files/upsdrv.rc6
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/www-apps/nut/files/upsdrv.rc6,v 1.1 2004/08/15 15:19:10 stuart Exp $

depend() {
	need net
	before upsd
}

start() {
	ebegin "Starting UPS drivers"
	/sbin/upsdrvctl start
	eend $? "Failed to start UPS drivers!"
}

stop() {
	ebegin "Stopping UPS drivers"
	/sbin/upsdrvctl stop
	eend $? "Failed to stop UPS drivers!"
}
