# net-misc/nut/files/upsd-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/nut/files/upsd-init,v 1.4 2004/07/15 00:06:11 agriffis Exp $

depend() {
	need net
	before upsmon
}

start() {
	ebegin "Starting upsd"
	/sbin/upsdrvctl start
	start-stop-daemon --start --quiet --startas /usr/sbin/upsd \
		--pidfile=/var/state/nut/upsd.pid
	eend $?
}

stop() {
	ebegin "Stopping upsd"
	start-stop-daemon --stop --quiet --pidfile=/var/state/nut/upsd.pid
	rc=$?
	/sbin/upsdrvctl stop
	eend $rc
}
