# net-misc/nut/files/upsdrv.rc6
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/nut/files/upsdrv.rc6,v 1.5 2004/07/15 00:06:11 agriffis Exp $

depend() {
	need net
	before upsd
}

start() {
	ebegin "Starting UPS drivers"
	/sbin/upsdrvctl start
	eend $? "Failed to start UPS drivers!"
}

stop() {
	ebegin "Stopping UPS drivers"
	/sbin/upsdrvctl stop
	eend $? "Failed to stop UPS drivers!"
}
