# net-misc/nut/files/upsmon.rc6
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/nut/files/upsmon.rc6,v 1.5 2004/04/01 05:56:18 robbat2 Exp $

depend() {
	need net
}

start() {
	ebegin "Starting upsmon"
	start-stop-daemon --start --quiet --exec /usr/sbin/upsmon
	eend $?
}

stop() {
	ebegin "Stopping upsmon"
	start-stop-daemon --stop --quiet --pidfile /var/run/upsmon.pid
	eend $?
}
