# www-servers/cherokee/files/cherokee-initd-1.2.0
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

opts="graceful reload"

depend() {
	need net
}

start() {
	ebegin "Starting Cherokee Web Server"
	start-stop-daemon --start -b \
		--exec /usr/sbin/cherokee
	eend $?
}

stop() {
	ebegin "Stopping Cherokee Web Server"
	start-stop-daemon --stop \
		--exec /usr/sbin/cherokee
	eend $?
}

graceful() {
	ebegin "Gracefully restarting Cherokee Web Server"
	start-stop-daemon \
		--exec /usr/sbin/cherokee \
		--signal HUP
	eend $?
}

reload() {
	ebegin "Restarting Cherokee Web Server and closing all open connections"
	start-stop-daemon \
		--exec /usr/sbin/cherokee \
		--signal USR1
	eend $?
}
