# sys-apps/usbd/files/usbd
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-apps/usbd/files/usbd,v 1.5 2005/08/18 00:25:10 vapier Exp $

depend() {
	need localmount
}

start() {
	ebegin "Starting usb daemon"
	/usr/sbin/usbd
	eend $?
}

stop() {
	ebegin "Stopping usb daemon"
	killall usbd
	eend $?
}
