# sys-apps/usbd/files/usbd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: 

depend() {
	need localmount
}

start() {
	ebegin "Starting usb daemon"
	/usr/sbin/usbd
	eend $?
}

stop() {
	ebegin "Stopping usb daemon"
	killall /usr/sbin/usbd
	eend $?
}
