# www-servers/plb/files/plb.rc6
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net
}

start() {
	ebegin "Starting plb"
	start-stop-daemon --start --quiet --exec /usr/bin/plb -- --daemonize --config=/etc/plb.conf
	eend $?
}

stop() {
	ebegin "Stopping plb"
	start-stop-daemon --stop --quiet --pidfile /var/run/plb.pid
	eend $?
}
