# sys-apps/hal/files/hald
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/sys-apps/hal/files/hald,v 1.2 2004/09/20 22:58:00 foser Exp $

depend() {

	need dbus hotplug

}

start() {

	ebegin "Starting Hardware Abstraction Layer daemon"

	start-stop-daemon --start -q --exec /usr/sbin/hald
	
	eend $?

}

stop() {

	ebegin "Stopping Hardware Abstraction Layer daemon"

	start-stop-daemon --stop -q --pidfile /var/run/hald/hald.pid 

	eend $?

}
