# www-servers/jboss/files/init.d/jboss
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/www-servers/jboss/files/init.d/jboss,v 1.1 2004/08/08 17:14:25 stuart Exp $

depend() {
	need net
}

start() {
	ebegin "Starting JBoss"
	cd $JBOSS_HOME/bin
	if [ -z "$SUBIT" ]; then
		eval $CMD_START >${JBOSS_CONSOLE} 2>&1 &
	else
		$SUBIT "$CMD_START >${JBOSS_CONSOLE} 2>&1 &" 
	fi
	eend $?
}

stop() {
	ebegin "Stopping JBoss"
	if [ -z "$SUBIT" ]; then
		$CMD_STOP
	else
		$SUBIT "$CMD_STOP"
	fi 
	eend $?
}
