# app-misc/tpconfig/files/tpconfig
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/app-misc/tpconfig/files/tpconfig,v 1.3 2003/04/29 20:07:56 mholzer Exp $

depend() {
        before xdm
}

start() {
	ebegin "Loading tpconfig options..."
	/usr/bin/tpconfig ${TPOPTIONS} > /dev/null 2> /tmp/tpconfig.start
	if [ -z "`cat /tmp/tpconfig.start 2> /dev/null`" ]; then
		eend 0
	else
		eend 1
	fi
}

stop() {
	ebegin "Stoping tpconfig (this does nothing)..."
	eend 0
}
