# dev-db/ingres/files/ingres-initd
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-db/ingres/files/ingres-initd,v 1.1 2009/10/04 00:08:23 patrick Exp $

depend() {
	use net
}

II_INSTALLATION=${SVCNAME:6:2}
export II_SYSTEM="/var/lib/ingres/Ingres${II_INSTALLATION}"

checkIngres() {

        if [ ! -d "${II_SYSTEM}" ]; then
		eerror "Ingres instance ${II_INSTALLATION} not found in ${II_SYSTEM}"
		return 1
	fi
}

start() {

	checkIngres || return 1 

	ebegin "Starting Ingres Instance ${II_INSTALLATION}"

	su ingres -c "
		. ${II_SYSTEM}/.ing${II_INSTALLATION}bash
		ingstart > /dev/null"
	RET=$?
	eend $RET

	if [ ! $RET ]; then
		eerror "Startup failed"
		eerror "You may find information in ${II_SYSTEM}/ingres/files/errlog.log"
		return 1
	fi

}

stop() {
	ebegin "Stopping Ingres Instance ${II_INSTALLATION}"

	checkIngres || return 1

        su ingres -c "
		. ${II_SYSTEM}/.ing${II_INSTALLATION}bash
		ingstop > /dev/null"
	
	RET=$?
	eend $RET

        if [ ! $RET ]; then  
                eerror "Stopping failed"
                eerror "You may try to run ingstop -f as user ingres"
		return 1
        fi
	
}

