# net-misc/bird/files/initrd-bird-1.0.11
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/bird/files/initrd-bird-1.0.11,v 1.1 2006/07/10 09:20:42 chainsaw Exp $

SOCK=/var/lib/run/bird.ctl

depend() {
	need net
	use logger
}

checkconfig() {
	if [ ! -f /etc/bird.conf ]; then
		eerror "Please create /etc/bird.conf"
		return 1
	fi
	return 0
}

start() {
	checkconfig || return $?

	ebegin "Starting BIRD"
	/usr/sbin/bird -c /etc/bird.conf -s "${SOCK}"
	eend $? "Failed to start BIRD"
}

stop() {
	ebegin "Stopping BIRD"
	if [ -f "${SOCK}" ]; then
		echo "down" | /usr/sbin/birdc &>/dev/null
		eend $? "Failed to stop BIRD"
	else
		eend 0
	fi
}

reload() {
	ebegin "Reloading BIRD configuration"

	eend $? "Failed to reload BIRD configuration"
}
