# sys-freebsd/freebsd-baselayout/files/dummy.initd
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

start() {
	# This is a fake, you know that?
	einfo "Fake"
}
