# app-admin/chef-solr/files/initd/chef-solr
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net rabbitmq
	after couchdb
}

CHEF_SOLR_BINARY=${CHEF_SOLR_BINARY:-/usr/bin/chef-solr}
CHEF_SOLR_PIDFILE=${CHEF_SOLR_PIDFILE:-/var/run/chef/solr.pid}
CHEF_SOLR_LOGFILE=${CHEF_SOLR_LOGFILE:-/var/log/chef/solr.log}
CONFIGFILE=${CONFIGFILE:-/etc/chef/solr.rb}
CHEF_USER=${CHEF_USER:-chef}
CHEF_GROUP=${CHEF_GROUP:-chef}

CHEF_DIRS=( /etc/chef /var/lib/chef /var/log/chef /var/run/chef )

fix_owners() {
	for i in ${CHEF_DIRS[@]} ${CONFIGFILE}; do
		chown ${CHEF_USER}:${CHEF_GROUP} ${i}
	done
	chmod 0600 ${CONFIGFILE}
}

start() {
	ebegin "Starting Chef SOLR"
	fix_owners
	start-stop-daemon --start --quiet --chuid ${CHEF_USER}:${CHEF_GROUP} \
		--pidfile ${CHEF_SOLR_PIDFILE} --name java --exec \
		${CHEF_SOLR_BINARY} -- -d -P ${CHEF_SOLR_PIDFILE} -c ${CONFIGFILE} \
		-L ${CHEF_SOLR_LOGFILE} ${CHEF_SOLR_OPTS} > /dev/null
	eend $?
}

stop() {
	ebegin "Shutting down Chef SOLR"
	start-stop-daemon --stop --quiet --pidfile ${CHEF_SOLR_PIDFILE}
	eend $?
}
