# net-analyzer/neti/files/neti-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-analyzer/neti/files/neti-init,v 1.3 2005/04/29 22:47:24 dragonheart Exp $

depend() {
    need net
}

start() {
	ebegin "Starting netiathome"
	start-stop-daemon --start --background --quiet --exec /usr/sbin/NETIStart
	eend $?
}

stop() {
	ebegin "Stopping netiathome"
	start-stop-daemon --stop --quiet --exec /usr/bin/tethereal
	eend $?
}
