# net-analyzer/neti/files/neti-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-analyzer/neti/files/neti-init,v 1.2 2004/07/14 23:05:29 agriffis Exp $

depend() {
    need net
}

start() {
	ebegin "Starting netiathome"
	start-stop-daemon --start --quiet --exec \
		/usr/bin/tethereal -iany -l -p -tad -V -R"udp||tcp||icmp||igmp" | neti
	eend $?
}

stop() {
	ebegin "Stopping netiathome"
	start-stop-daemon --stop --quiet --exec /usr/bin/tethereal
	eend $?
}
