# sys-cluster/ipvsadm/files/ipvsadm-init
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-cluster/ipvsadm/files/ipvsadm-init,v 1.3 2003/12/27 16:47:17 iggy Exp $

opts="start stop save"

depend() {
        need net
}

checkconfig() {
	if [ ! -f /var/lib/ipvsadm/rules-save ] ; then
		eerror "Not starting ipvsadm. First create some rules then run"
		eerror "/etc/init.d/ipvsadm save"
		return 1
	fi
}

start() {
	ebegin "Starting ipvsadm"
	einfo "Clearing previous rules"
	ipvsadm -C
	ipvsadm-restore < /var/lib/ipvsadm/rules-save
	eend $?
}

stop() {
	ebegin "Stopping ipvsadm"
	eend $?
}

save() {
	ebegin "Saving ipvsadm state"
	ipvsadm-save > /var/lib/ipvsadm/rules-save
	eend $?
}
