# dev-db/postgresql/files/7.1.3/postgresql
#!/sbin/runscript

PREFIX="/usr"
PGDATA="/var/db/postgresql/data"
PGUSER="postgres"
PATH="/bin:/usr/bin"
DAEMON="$PREFIX/bin/pg_ctl"

opts="start stop restart status"

depend() {
    need net
}

start() {
    ebegin "Starting PostgreSQL"
    su - $PGUSER -c "$DAEMON start -D '$PGDATA' -s -o -i" | logger 1>&2 &
    eend $?
}

stop() {
    ebegin "Stopping PostgreSQL"
    su - $PGUSER -c "$DAEMON stop -D '$PGDATA' -s -m fast"
    eend $?
}

restart() {
    ebegin "Restarting PostgreSQL"
    su - $PGUSER -c "$DAEMON restart -D '$PGDATA' -s -m fast" &
    eend $?
}

status() {
    su - $PGUSER -c "$DAEMON status -D '$PGDATA'"
}
