# net-analyzer/darkstat/files/darkstat-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2 
# $Header:

depend() {
        need net
}

start() {
	ebegin "Starting darkstat"
	start-stop-daemon --start --quiet --exec /usr/sbin/darkstat -- --detach -d /var/spool/darkstat ${DARKSTAT_OPTS} >/dev/null
	eend $?
}

stop() {
	ebegin "Stopping darkstat"
	start-stop-daemon --stop --quiet --retry 5 --exec /usr/sbin/darkstat
	eend $?
}
