# net-analyzer/darkstat/files/darkstat-init2
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2 
# $Header: /var/cvsroot/gentoo-x86/net-analyzer/darkstat/files/darkstat-init2,v 1.1 2006/07/08 02:30:32 vanquirius Exp $

depend() {
        need net
}

start() {
	ebegin "Starting darkstat"
	start-stop-daemon --start --quiet --exec /usr/sbin/darkstat -- -d /var/spool/darkstat ${DARKSTAT_OPTS} >/dev/null
	eend $?
}

stop() {
	ebegin "Stopping darkstat"
	start-stop-daemon --stop --quiet --retry 5 --exec /usr/sbin/darkstat
	eend $?
}
