# net-im/jive-messenger/files/init.d/jive
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting Jive Messenger"
	MESSENGER_HOME=/opt/jive_messenger
	MESSENGER_LIB="${MESSENGER_HOME}/lib"
	MESSENGER_OPTS="-Xmx256m -DmessengerHome=${MESSENGER_HOME} -Dmessenger.lib.dir=${MESSENGER_LIB}"
	sudo -u jive $(java-config --java) -server ${MESSENGER_OPTS} -classpath "${MESSENGER_LIB}/startup.jar" -jar "${MESSENGER_LIB}/startup.jar" > /var/log/jive.log 2>&1 &
	echo $! > /var/run/jive.pid
	JIVE_PID=$(</var/run/jive.pid)
	if [ -n "${JIVE_PID}" ]; then 
		ps ax|grep "^ *$(< /var/run/jive.pid)" > /dev/null
		eend $?
	else 
		eend 1
	fi
}

stop() {
	ebegin "Stopping Jive Messenger"
	kill $(< /var/run/jive.pid)
	sleep 1
	if ps ax|grep "^ *$(< /var/run/jive.pid)" > /dev/null; then
		eend 1
	else
		rm /var/run/jive.pid
		eend 0
	fi
}
