# net-misc/dhcpcd/files/dhcpcd.service
[Unit]
Description=Lightweight DHCP client daemon

[Service]
ExecStart=/sbin/dhcpcd -q --nobackground

[Install]
WantedBy=network.target
