# net-dns/pdns-recursor/files/precursor
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-dns/pdns-recursor/files/precursor,v 1.1 2006/04/30 00:14:42 swegener Exp $

depend() {
	need net
}

start() {
	ebegin "Starting PowerDNS Recursor"
	/usr/sbin/pdns_recursor --daemon &>/dev/null
	eend $?
}

stop() {
	ebegin "Stopping PowerDNS Recursor"
	/usr/sbin/rec_control quit &>/dev/null
	eend $?
}
