# net-dns/pdns-recursor/files/precursor
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-dns/pdns-recursor/files/precursor,v 1.3 2006/10/27 16:40:31 swegener Exp $

depend() {
	need net
}

start() {
	ebegin "Starting PowerDNS Recursor"
	/usr/sbin/pdns_recursor --daemon=yes &>/dev/null
	eend $?
}

stop() {
	ebegin "Stopping PowerDNS Recursor"
	start-stop-daemon --stop --quiet --pidfile=/var/run/pdns_recursor.pid
	eend $?
}
