# net-mail/courier/files/courier-imapd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-mail/courier/files/courier-imapd,v 1.4 2004/03/06 03:15:17 vapier Exp $

depend() {
	need net courier-authdaemond
}

start() {
	ebegin "Starting courier imapd"
	/usr/share/courier/imapd start
	eend $?
}

stop() {
	ebegin "Stopping courier imapd"
	/usr/share/courier/imapd stop
	eend $?
}
