# net-mail/courier/files/courier-pop3d
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-mail/courier/files/courier-pop3d,v 1.4 2004/03/06 03:15:17 vapier Exp $

depend() {
	need net courier-authdaemond
}

start() {
	ebegin "Starting courier pop3d"
	/usr/share/courier/pop3d start
	eend $?
}

stop() {
	ebegin "Stopping courier pop3d"
	/usr/share/courier/pop3d stop
	eend $?
}
