# net-mail/courier/files/courier
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# /space/gentoo/cvsroot/gentoo-x86/net-misc/bind/files/named.rc6,v 1.3 2002/02/09 01:08:04 woodchip Exp

depend() {
	need net courier-authdaemond courier-mta courier-esmtpd courier-esmtpd-ssl courier-esmtpd-msa courier-imapd courier-imapd-ssl courier-pop3d courier-pop3d-ssl
}

start() {
	ebegin "Courier MTA startet"
	eend $?
}

stop() {
	ebegin "Stopping Courier MTA"
	/etc/init.d/courier-authdaemond stop
	eend $?
}
