# net-mail/courier/files/courier-mta
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-mail/courier/files/courier-mta,v 1.5 2004/03/06 03:15:17 vapier Exp $

depend() {
	need net courier-authdaemond courier-filterd
}

start() {
	ebegin "Starting courier mta"
	/usr/sbin/courier start
	eend $?
}

stop() {
	ebegin "Stopping courier mta"
	/usr/sbin/courier stop
	eend $?
}
