# net-mail/courier/files/courier
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-mail/courier/files/courier,v 1.5 2004/03/06 03:15:17 vapier Exp $

DAEMONLIST="courier-authdaemond courier-filterd courier-mta courier-esmtpd courier-esmtpd-ssl courier-imapd courier-imapd-ssl courier-pop3d courier-pop3d-ssl"

depend() {
	need net
}

start() {
	ebegin "Starting Courier MTA"
	for d in $DAEMONLIST
	do
		[ -f /mnt/.init.d/started/$d ] || /etc/init.d/$d start || true
	done
	eend $?
}

stop() {
	ebegin "Stopping Courier MTA"
	for d in $DAEMONLIST
	do
		[ -f /mnt/.init.d/started/$d ] && /etc/init.d/$d stop || true
	done
	eend $?
}
