# net-mail/courier/files/courier-mta
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# /space/gentoo/cvsroot/gentoo-x86/net-mail/courier/files/courier-mta

depend() {
	need net courier-authdaemond
}

start() {
	ebegin "Starting courier-mta"
	/usr/sbin/courier start
	eend $?
}

stop() {
	ebegin "Stopping courier-mta"
	/usr/sbin/courier stop
	eend $?
}
