# net-mail/courier/files/courier-filterd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-mail/courier/files/courier-filterd,v 1.3 2004/03/06 03:15:17 vapier Exp $

depend() {
	need net
}

start() {
	ebegin "Starting courier filterd"
	/usr/sbin/courierfilter start
	eend $?
}

stop() {
	ebegin "Stopping courier filterd"
	/usr/sbin/courierfilter stop
	eend $?
}

restart() {
	ebegin "Restarting courier filterd"
	/usr/sbin/courierfilter restart
	eend $?
}
