# net-mail/courier/files/courier-esmtpd-msa
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-mail/courier/files/courier-esmtpd-msa,v 1.4 2004/03/06 03:15:17 vapier Exp $

depend() {
	need net courier-authdaemond
}

start() {
	ebegin "Starting courier esmtpd-msa"
	/usr/sbin/esmtpd-msa start
	eend $?
}

stop() {
	ebegin "Stopping courier esmtpd-msa"
	/usr/sbin/esmtpd-msa stop
	eend $?
}
