# sys-fs/gfs/files/gfs
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/sys-fs/gfs/files/gfs,v 1.2 2005/03/19 22:04:36 xmerlin Exp $

depend() {
	provide cluster
	need net
	need ccs fenced
	after cman
	use dns logger
}

checkconfig() {
	if [ ! -f /etc/ntp.conf ] ; then
		eerror "Please create /etc/ntp.conf"
		eerror "Sample conf: /usr/share/ntp/ntp.conf"
		return 1
	fi
	return 0
}

start() {
	checkconfig || return $?

	einfo "Starting gfs cluster:"
        ebegin "Loading needed kernel modules"

        if [ ! -f /proc/cluster/lock_dlm ]; then
                modprobe lock_dlm
        fi
	if [ ! -f /proc/fs/gfs ]; then
		modprobe gfs
	fi
	eend $? "Failed to load needed kernel modules"

}

stop() {
	einfo "Stopping gfs cluster:"
	ebegin "Removing gfs module"
	modprobe -r gfs &>/dev/null
	eend $?
	
}
