# sys-fs/gfs/files/gfs-mount
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/sys-fs/gfs/files/gfs-mount,v 1.1 2005/03/19 17:32:19 xmerlin Exp $

depend() {
	need net 
	after gnbd-client gfs
	use dns logger
}

start() {
	GFS=`cat /etc/fstab | grep 'gfs'`
	if [ -n "$GFS" ] ; then
		ebegin "Mounting all gfs filesystems"
		cat /etc/fstab | awk '/[\t\s]gfs[\t\s]/ { print $1,$2 }' | xargs mount -v -t gfs
		eend $? "Failed to mount all gfs filesystems"
	fi
}

stop() {
	GFS=`cat /etc/fstab | grep 'gfs'`
	if [ -n "$GFS" ] ; then
		ebegin "Unmounting all gfs filesystems"
		cat /etc/fstab | awk '/[\t\s]gfs[\t\s]/ { print $2 }' | xargs umount -v -t gfs
		eend $? "Failed to umount all gfs filesystems"
	fi
}
