# net-dialup/dwun/files/dwun
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-dialup/dwun/files/dwun,v 1.2 2004/03/04 19:10:18 vapier Exp $

depend() {
	need net
}

DAEMON=/usr/sbin/dwun
PIDFILE=/var/run/dwun.pid

start() {
	ebegin "Starting dwun"
	start-stop-daemon --start --pidfile ${PIDFILE} --startas ${DAEMON}
	eend $?
}

stop() {
	ebegin "Stopping dwun"
	start-stop-daemon --stop --signal 1 --pidfile ${PIDFILE}
	eend $?
}
